# question_generation/grouped_attackseqs/183.json
{
    "file_name": "Gallmaker New Attack Group Eschews Malware to Live off the Land",
    "tactic_label": true,
    "rewrite": {
        "Initial Access": "The Gallmaker group initiated their attack by delivering malicious Office documents, likely through spear-phishing emails, with titles related to government, military, and diplomatic themes. These documents, with filenames in English or Cyrillic, attempted to exploit the Microsoft Office Dynamic Data Exchange (DDE) protocol to execute commands in memory on the victim's system when the content was enabled, thus avoiding leaving artifacts on disk.",
        "Execution": "Upon gaining access to a victim's device, Gallmaker executed various tools including WindowsRoamingToolsTask for scheduling PowerShell scripts and tasks, a 'reverse_tcp' payload from Metasploit via obfuscated PowerShell, a legitimate version of the WinZip console for command execution and potential data exfiltration, and the Rex PowerShell library from GitHub for creating and manipulating PowerShell scripts for Metasploit exploits.",
        "Defense Evasion": "Gallmaker's use of living off the land (LotL) tactics and publicly available hack tools, as well as their practice of deleting tools from victim machines after use, helped them avoid detection. They exploited the DDE protocol to execute code remotely without leaving artifacts on disk.",
        "Command and Control": "Gallmaker used three primary IP addresses for command-and-control (C&C) infrastructure to communicate with infected devices. They also likely used the WinZip console to archive data for exfiltration.",
        "Exfiltration": "The legitimate version of the WinZip console was likely used by Gallmaker to archive data, probably for the purpose of exfiltration.",
        "Others": "Gallmaker targeted government, military, and defense sectors, including overseas embassies of an Eastern European country and military and defense targets in the Middle East. The group has been active since at least December 2017, with the most recent activity observed in June 2018. Their activities are indicative of a cyber espionage campaign, likely state-sponsored. Symantec's Targeted Attack Analytics (TAA) technology played a crucial role in detecting Gallmaker's activities, which relied on LotL tactics and publicly available hack tools to avoid detection. Protections against Gallmaker attacks include System Infected: Meterpreter Reverse TCP and W97M.Downloader. Network protection products also detect activity associated with Gallmaker. Indicators of Compromise (IoCs) specific to Gallmaker have been updated, and Yara rules and methods used by Gallmaker are shared for detection of unauthorized activity."
    },
    "technique_label": true,
    "triplet_groups": {
        "Initial Access": {
            "T1566-Phishing": [
                {
                    "Subject": "Gallmaker group",
                    "SubjectType": "threat-actor",
                    "Relation": "deliver",
                    "Object": "malicious Office documents",
                    "ObjectType": "file",
                    "tactic": "Initial Access",
                    "technique": [
                        "T1566-Phishing"
                    ],
                    "id": 1
                },
                {
                    "Subject": "malicious Office documents",
                    "SubjectType": "file",
                    "Relation": "exploit",
                    "Object": "Microsoft Office Dynamic Data Exchange (DDE) protocol",
                    "ObjectType": "vulnerability",
                    "tactic": "Initial Access",
                    "technique": [
                        "T1566-Phishing"
                    ],
                    "id": 2
                }
            ]
        },
        "Execution": {
            "T1059-Command and Scripting Interpreter": [
                {
                    "Subject": "Gallmaker",
                    "SubjectType": "threat-actor",
                    "Relation": "execute",
                    "Object": "WindowsRoamingToolsTask",
                    "ObjectType": "tool",
                    "tactic": "Execution",
                    "technique": [
                        "T1059-Command and Scripting Interpreter",
                        "T1053-Scheduled Task/Job"
                    ],
                    "id": 3
                },
                {
                    "Subject": "Gallmaker",
                    "SubjectType": "threat-actor",
                    "Relation": "use",
                    "Object": "reverse_tcp payload",
                    "ObjectType": "tool",
                    "tactic": "Execution",
                    "technique": [
                        "T1059-Command and Scripting Interpreter",
                        "T1203-Exploitation for Client Execution"
                    ],
                    "id": 4
                },
                {
                    "Subject": "Gallmaker",
                    "SubjectType": "threat-actor",
                    "Relation": "use",
                    "Object": "WinZip console",
                    "ObjectType": "tool",
                    "tactic": "Execution",
                    "technique": [
                        "T1059-Command and Scripting Interpreter"
                    ],
                    "id": 5
                },
                {
                    "Subject": "Gallmaker",
                    "SubjectType": "threat-actor",
                    "Relation": "use",
                    "Object": "Rex PowerShell library",
                    "ObjectType": "tool",
                    "tactic": "Execution",
                    "technique": [
                        "T1059-Command and Scripting Interpreter",
                        "T1559-Inter-Process Communication"
                    ],
                    "id": 6
                }
            ],
            "T1053-Scheduled Task/Job": [
                {
                    "Subject": "Gallmaker",
                    "SubjectType": "threat-actor",
                    "Relation": "execute",
                    "Object": "WindowsRoamingToolsTask",
                    "ObjectType": "tool",
                    "tactic": "Execution",
                    "technique": [
                        "T1059-Command and Scripting Interpreter",
                        "T1053-Scheduled Task/Job"
                    ],
                    "id": 3
                }
            ],
            "T1203-Exploitation for Client Execution": [
                {
                    "Subject": "Gallmaker",
                    "SubjectType": "threat-actor",
                    "Relation": "use",
                    "Object": "reverse_tcp payload",
                    "ObjectType": "tool",
                    "tactic": "Execution",
                    "technique": [
                        "T1059-Command and Scripting Interpreter",
                        "T1203-Exploitation for Client Execution"
                    ],
                    "id": 4
                }
            ],
            "T1559-Inter-Process Communication": [
                {
                    "Subject": "Gallmaker",
                    "SubjectType": "threat-actor",
                    "Relation": "use",
                    "Object": "Rex PowerShell library",
                    "ObjectType": "tool",
                    "tactic": "Execution",
                    "technique": [
                        "T1059-Command and Scripting Interpreter",
                        "T1559-Inter-Process Communication"
                    ],
                    "id": 6
                }
            ]
        },
        "Defense Evasion": {
            "T1070-Indicator Removal": [
                {
                    "Subject": "Gallmaker",
                    "SubjectType": "threat-actor",
                    "Relation": "use",
                    "Object": "living off the land (LotL) tactics",
                    "ObjectType": "Tactic",
                    "tactic": "Defense Evasion",
                    "technique": [
                        "T1070-Indicator Removal"
                    ],
                    "id": 7
                },
                {
                    "Subject": "Gallmaker",
                    "SubjectType": "threat-actor",
                    "Relation": "exploit",
                    "Object": "DDE protocol",
                    "ObjectType": "vulnerability",
                    "tactic": "Defense Evasion",
                    "technique": [
                        "T1070-Indicator Removal"
                    ],
                    "id": 8
                }
            ]
        },
        "Command and Control": {
            "T1573-Encrypted Channel": [
                {
                    "Subject": "Gallmaker",
                    "SubjectType": "threat-actor",
                    "Relation": "use",
                    "Object": "IP addresses",
                    "ObjectType": "ipv4-addr",
                    "tactic": "Command and Control",
                    "technique": [
                        "T1573-Encrypted Channel"
                    ],
                    "id": 9
                }
            ],
            "T1105-Ingress Tool Transfer": [
                {
                    "Subject": "Gallmaker",
                    "SubjectType": "threat-actor",
                    "Relation": "use",
                    "Object": "WinZip console",
                    "ObjectType": "tool",
                    "tactic": "Command and Control",
                    "technique": [
                        "T1105-Ingress Tool Transfer"
                    ],
                    "id": 10
                }
            ]
        },
        "Exfiltration": {
            "T1052-Exfiltration Over Physical Medium": [
                {
                    "Subject": "Gallmaker",
                    "SubjectType": "threat-actor",
                    "Relation": "use",
                    "Object": "WinZip console",
                    "ObjectType": "tool",
                    "tactic": "Exfiltration",
                    "technique": [
                        "T1052-Exfiltration Over Physical Medium"
                    ],
                    "id": 11
                },
                {
                    "Subject": "Gallmaker",
                    "SubjectType": "threat-actor",
                    "Relation": "exfiltrate-to",
                    "Object": "data",
                    "ObjectType": "artifacts",
                    "tactic": "Exfiltration",
                    "technique": [
                        "T1052-Exfiltration Over Physical Medium"
                    ],
                    "id": 12
                }
            ]
        }
    }
}