# question_generation/grouped_attackseqs/22.json
{
    "file_name": "New Perl Botnet (Tuyul) Found with Possible Indonesian Attribution",
    "tactic_label": true,
    "rewrite": {
        "Initial Access": "On January 15, 2020, F5 threat researchers detected a new campaign exploiting a vulnerability in PHPUnit systems (CVE-2017-9841) to install an IRC bot named Tuyul. The malware specifically targets unpatched PHPUnit, a PHP unit testing framework, allowing the attacker to inject arbitrary PHP code on the server.",
        "Execution": "The Tuyul botnet, written in Perl, is used to execute commands on victim systems. The malware source code includes functionalities such as promoting bots to higher privilege (pro), demoting bots (noob), updating malware, connecting to reverse shells (bc), executing shell commands (shell), running Perl scripts (perl), and terminating the bot's connection (terminate, suicide).",
        "Persistence": "To maintain persistence, the Tuyul malware periodically downloads a bash script named 'cron' to ensure its continued presence on the system. It also includes functionality to detect and terminate processes from other malware infections.",
        "Defense Evasion": "The Tuyul malware employs a binary version that goes undetected by antivirus engines, allowing it to operate on the machine without signature-based detection. Additionally, the malware uses various PHP web shells to gain access to servers, including an open-source application called Tiny File Manager and an unknown advanced web shell.",
        "Command and Control": "The Tuyul botnet uses the domain irc[.]localroot[.]xyz as its command and control (C&C) server, which hosts the malware files and has a valid SSL certificate. The attacker also maintains two development servers for script development and testing, and for logging successful infections. The C&C server and development servers are protected by Cloudflare.",
        "Others": "The Tuyul botnet has possible Indonesian attribution based on several clues, including the time zone set to Asia/Jakarta, the botnet name, admin nicknames, and repository server. The botnet's peak size was 366 victim systems, which is smaller compared to other botnets. The infected servers are primarily hosted on cloud services such as AWS, DigitalOcean, and Linode, with a significant number located in the U.S. and U.K. The botnet's name, Tuyul, refers to a figure in Southeast Asian folklore, particularly Indonesian, which is an undead infant invoked for theft and other crimes. The botnet is actively being worked on and is expected to grow."
    },
    "technique_label": true,
    "triplet_groups": {
        "Initial Access": {
            "T1190-Exploit Public-Facing Application": [
                {
                    "Subject": "F5",
                    "SubjectType": "CTI vender",
                    "Relation": "detected",
                    "Object": "campaign",
                    "ObjectType": "campaign",
                    "tactic": "Initial Access",
                    "technique": [
                        "T1190-Exploit Public-Facing Application"
                    ],
                    "id": 1
                },
                {
                    "Subject": "campaign",
                    "SubjectType": "campaign",
                    "Relation": "exploit",
                    "Object": "CVE-2017-9841",
                    "ObjectType": "vulnerability",
                    "tactic": "Initial Access",
                    "technique": [
                        "T1190-Exploit Public-Facing Application"
                    ],
                    "id": 2
                },
                {
                    "Subject": "Tuyul",
                    "SubjectType": "Malware",
                    "Relation": "target",
                    "Object": "PHPUnit",
                    "ObjectType": "software",
                    "tactic": "Initial Access",
                    "technique": [
                        "T1190-Exploit Public-Facing Application"
                    ],
                    "id": 3
                }
            ]
        },
        "Execution": {
            "T1059-Command and Scripting Interpreter": [
                {
                    "Subject": "The Tuyul botnet",
                    "SubjectType": "Malware",
                    "Relation": "use",
                    "Object": "Perl",
                    "ObjectType": "software",
                    "tactic": "Execution",
                    "technique": [
                        "T1059-Command and Scripting Interpreter"
                    ],
                    "id": 4
                },
                {
                    "Subject": "The Tuyul botnet",
                    "SubjectType": "Malware",
                    "Relation": "execute",
                    "Object": "commands",
                    "ObjectType": "artifact",
                    "tactic": "Execution",
                    "technique": [
                        "T1059-Command and Scripting Interpreter"
                    ],
                    "id": 5
                }
            ]
        },
        "Persistence": {
            "T1053-Scheduled Task/Job": [
                {
                    "Subject": "Tuyul malware",
                    "SubjectType": "Malware",
                    "Relation": "download",
                    "Object": "cron",
                    "ObjectType": "file",
                    "tactic": "Persistence",
                    "technique": [
                        "T1053-Scheduled Task/Job"
                    ],
                    "id": 6
                },
                {
                    "Subject": "Tuyul malware",
                    "SubjectType": "Malware",
                    "Relation": "include",
                    "Object": "functionality",
                    "ObjectType": "feature",
                    "tactic": "Persistence",
                    "technique": [
                        "T1053-Scheduled Task/Job"
                    ],
                    "id": 7
                }
            ]
        },
        "Defense Evasion": {
            "T1564-Hide Artifacts": [
                {
                    "Subject": "The Tuyul malware",
                    "SubjectType": "Malware",
                    "Relation": "use",
                    "Object": "binary version",
                    "ObjectType": "artifact",
                    "tactic": "Defense Evasion",
                    "technique": [
                        "T1564-Hide Artifacts"
                    ],
                    "id": 8
                }
            ]
        },
        "Command and Control": {
            "T1090-Proxy": [
                {
                    "Subject": "Tuyul botnet",
                    "SubjectType": "Malware",
                    "Relation": "use",
                    "Object": "irc[.]localroot[.]xyz",
                    "ObjectType": "domain",
                    "tactic": "Command and Control",
                    "technique": [
                        "T1090-Proxy"
                    ],
                    "id": 12
                },
                {
                    "Subject": "attacker",
                    "SubjectType": "threat-actor",
                    "Relation": "maintain",
                    "Object": "development servers",
                    "ObjectType": "infrastrcucture",
                    "tactic": "Command and Control",
                    "technique": [
                        "T1090-Proxy"
                    ],
                    "id": 14
                }
            ],
            "T1071-Application Layer Protocol": [
                {
                    "Subject": "irc[.]localroot[.]xyz",
                    "SubjectType": "domain",
                    "Relation": "host",
                    "Object": "malware files",
                    "ObjectType": "file",
                    "tactic": "Command and Control",
                    "technique": [
                        "T1071-Application Layer Protocol"
                    ],
                    "id": 13
                }
            ]
        }
    }
}