# question_generation/grouped_attackseqs/384.json
{
    "file_name": "DNSpionage Campaign Targets Middle East",
    "tactic_label": true,
    "rewrite": {
        "Reconnaissance": "The DNSpionage campaign involved meticulous reconnaissance, as evidenced by the attackers' understanding of the victims' network infrastructure to remain undetected. They targeted Lebanon and the UAE, focusing on .gov domains and a private Lebanese airline. The attackers likely used spear-phishing via email or social media platforms like LinkedIn to distribute malicious Microsoft Office documents with embedded macros, disguised as job postings on fake websites such as hr-wipro[.]com and hr-suncor[.]com.",
        "Resource Development": "The DNSpionage campaign demonstrated resource development by the attackers, who created two fake websites with job postings to compromise targets. They also generated Let's Encrypt certificates for redirected domains during DNS compromises, although it is unclear if these redirections were successful.",
        "Initial Access": "Initial access in the DNSpionage campaign was attempted through malicious websites hosting job listings, which redirected to legitimate company sites. The attackers likely distributed malicious Microsoft Office documents via spear-phishing to compromise targets.",
        "Execution": "The DNSpionage malware executes in two steps: first, it decodes a PE file encoded with base64 when the document is opened, and second, it renames and executes the file as a scheduled task when the document is closed. This process requires human interaction and is designed to avoid sandbox detection.",
        "Persistence": "The DNSpionage malware maintains persistence by creating a scheduled task named 'chromiumupdater v 37.5.0' to execute the dropped binary every minute.",
        "Defense Evasion": "The DNSpionage malware employs defense evasion techniques by using password-protected macros in Microsoft Word to prevent victims from accessing the macro code and by using classical string obfuscation to avoid detection.",
        "Discovery": "The DNSpionage malware's initial reconnaissance phase involves commands to retrieve the infected user's username, hostname, and domain information, indicating a discovery process to understand the environment.",
        "Command and Control": "The DNSpionage malware uses HTTP and DNS for command and control communication with the C2 server. It crafts DNS requests to register the infected system and receive further instructions, and uses HTTP requests to retrieve configurations and orders from a fake Wikipedia page.",
        "Impact": "The DNS redirection part of the DNSpionage campaign could have a significant impact, as the attackers were able to intercept traffic destined for the redirected hostnames, potentially harvesting email and VPN credentials, and possibly enabling further attacks or blackmail.",
        "Others": "The DNSpionage campaign was discovered by Cisco Talos, authored by Warren Mercer and Paul Rascagneres. The campaign targeted Lebanon and the UAE, with no clear connection to other campaigns or actors. The attackers used Windows malware and DNS exfiltration and redirection, with five attacks launched this year, including one recently. The campaign highlights the need for strong endpoint and network protection against advanced actors with significant targets."
    },
    "technique_label": true,
    "triplet_groups": {
        "Initial Access": {
            "T1566-Phishing": [
                {
                    "Subject": "DNSpionage campaign",
                    "SubjectType": "campaign",
                    "Relation": "attempt",
                    "Object": "initial access",
                    "ObjectType": "attack-pattern",
                    "tactic": "Initial Access",
                    "technique": [
                        "T1566-Phishing"
                    ],
                    "id": 1
                },
                {
                    "Subject": "attackers",
                    "SubjectType": "threat-actor",
                    "Relation": "distribute",
                    "Object": "malicious Microsoft Office documents",
                    "ObjectType": "artifact",
                    "tactic": "Initial Access",
                    "technique": [
                        "T1566-Phishing"
                    ],
                    "id": 2
                }
            ]
        },
        "Execution": {
            "T1053-Scheduled Task/Job": [
                {
                    "Subject": "DNSpionage malware",
                    "SubjectType": "Malware",
                    "Relation": "execute",
                    "Object": "PE file",
                    "ObjectType": "file",
                    "tactic": "Execution",
                    "technique": [
                        "T1053-Scheduled Task/Job"
                    ],
                    "id": 3
                }
            ],
            "T1204-User Execution": [
                {
                    "Subject": "DNSpionage malware",
                    "SubjectType": "Malware",
                    "Relation": "require",
                    "Object": "human interaction",
                    "ObjectType": "artifact",
                    "tactic": "Execution",
                    "technique": [
                        "T1204-User Execution"
                    ],
                    "id": 4
                }
            ]
        },
        "Defense Evasion": {
            "T1027-Obfuscated Files or Information": [
                {
                    "Subject": "DNSpionage malware",
                    "SubjectType": "malware",
                    "Relation": "employ",
                    "Object": "defense evasion",
                    "ObjectType": "technique",
                    "tactic": "Defense Evasion",
                    "technique": [
                        "T1027-Obfuscated Files or Information"
                    ],
                    "id": 5
                },
                {
                    "Subject": "DNSpionage malware",
                    "SubjectType": "malware",
                    "Relation": "use",
                    "Object": "string obfuscation",
                    "ObjectType": "technique",
                    "tactic": "Defense Evasion",
                    "technique": [
                        "T1027-Obfuscated Files or Information"
                    ],
                    "id": 7
                }
            ]
        },
        "Command and Control": {
            "T1071-Application Layer Protocol": [
                {
                    "Subject": "DNSpionage malware",
                    "SubjectType": "Malware",
                    "Relation": "use",
                    "Object": "HTTP and DNS",
                    "ObjectType": "network traffic",
                    "tactic": "Command and Control",
                    "technique": [
                        "T1071-Application Layer Protocol"
                    ],
                    "id": 8
                },
                {
                    "Subject": "DNSpionage malware",
                    "SubjectType": "Malware",
                    "Relation": "communicate-with",
                    "Object": "C2 server",
                    "ObjectType": "infrastructure",
                    "tactic": "Command and Control",
                    "technique": [
                        "T1071-Application Layer Protocol"
                    ],
                    "id": 9
                }
            ]
        },
        "Persistence": {
            "T1053-Scheduled Task/Job": [
                {
                    "Subject": "DNSpionage malware",
                    "SubjectType": "Malware",
                    "Relation": "maintain",
                    "Object": "persistence",
                    "ObjectType": "Tactic",
                    "tactic": "Persistence",
                    "technique": [
                        "T1053-Scheduled Task/Job"
                    ],
                    "id": 10
                },
                {
                    "Subject": "DNSpionage malware",
                    "SubjectType": "Malware",
                    "Relation": "create",
                    "Object": "scheduled task",
                    "ObjectType": "Technique",
                    "tactic": "Persistence",
                    "technique": [
                        "T1053-Scheduled Task/Job"
                    ],
                    "id": 11
                },
                {
                    "Subject": "scheduled task",
                    "SubjectType": "Technique",
                    "Relation": "named",
                    "Object": "chromiumupdater v 37.5.0",
                    "ObjectType": "artifact",
                    "tactic": "Persistence",
                    "technique": [
                        "T1053-Scheduled Task/Job"
                    ],
                    "id": 12
                },
                {
                    "Subject": "scheduled task",
                    "SubjectType": "Technique",
                    "Relation": "execute",
                    "Object": "dropped binary",
                    "ObjectType": "file",
                    "tactic": "Persistence",
                    "technique": [
                        "T1053-Scheduled Task/Job"
                    ],
                    "id": 13
                }
            ]
        },
        "Discovery": {
            "T1087-Account Discovery": [
                {
                    "Subject": "DNSpionage malware",
                    "SubjectType": "Malware",
                    "Relation": "retrieve",
                    "Object": "username, hostname, and domain information",
                    "ObjectType": "artifacts",
                    "tactic": "Discovery",
                    "technique": [
                        "T1087-Account Discovery"
                    ],
                    "id": 14
                }
            ]
        },
        "Reconnaissance": {
            "T1590-Gather Victim Network Information": [
                {
                    "Subject": "The DNSpionage campaign",
                    "SubjectType": "campaign",
                    "Relation": "involve",
                    "Object": "reconnaissance",
                    "ObjectType": "Tactic",
                    "tactic": "Reconnaissance",
                    "technique": [
                        "T1590-Gather Victim Network Information"
                    ],
                    "id": 15
                },
                {
                    "Subject": "The attackers",
                    "SubjectType": "threat-actor",
                    "Relation": "understand",
                    "Object": "victims' network infrastructure",
                    "ObjectType": "infrastructure",
                    "tactic": "Reconnaissance",
                    "technique": [
                        "T1590-Gather Victim Network Information"
                    ],
                    "id": 16
                }
            ],
            "T1591-Gather Victim Org Information": [
                {
                    "Subject": "The attackers",
                    "SubjectType": "threat-actor",
                    "Relation": "focus-on",
                    "Object": ".gov domains and a private Lebanese airline",
                    "ObjectType": "identity",
                    "tactic": "Reconnaissance",
                    "technique": [
                        "T1591-Gather Victim Org Information"
                    ],
                    "id": 18
                }
            ],
            "T1598-Phishing for Information": [
                {
                    "Subject": "The attackers",
                    "SubjectType": "threat-actor",
                    "Relation": "use",
                    "Object": "spear-phishing",
                    "ObjectType": "attack-pattern",
                    "tactic": "Reconnaissance",
                    "technique": [
                        "T1598-Phishing for Information"
                    ],
                    "id": 19
                },
                {
                    "Subject": "spear-phishing",
                    "SubjectType": "attack-pattern",
                    "Relation": "distribute",
                    "Object": "malicious Microsoft Office documents",
                    "ObjectType": "artifact",
                    "tactic": "Reconnaissance",
                    "technique": [
                        "T1598-Phishing for Information"
                    ],
                    "id": 20
                },
                {
                    "Subject": "malicious Microsoft Office documents",
                    "SubjectType": "artifact",
                    "Relation": "embed",
                    "Object": "macros",
                    "ObjectType": "artifact",
                    "tactic": "Reconnaissance",
                    "technique": [
                        "T1598-Phishing for Information"
                    ],
                    "id": 21
                },
                {
                    "Subject": "malicious Microsoft Office documents",
                    "SubjectType": "artifact",
                    "Relation": "disguise-as",
                    "Object": "job postings",
                    "ObjectType": "Email Message",
                    "tactic": "Reconnaissance",
                    "technique": [
                        "T1598-Phishing for Information"
                    ],
                    "id": 22
                },
                {
                    "Subject": "job postings",
                    "SubjectType": "Email Message",
                    "Relation": "host-on",
                    "Object": "fake websites",
                    "ObjectType": "domain-name",
                    "tactic": "Reconnaissance",
                    "technique": [
                        "T1598-Phishing for Information"
                    ],
                    "id": 23
                },
                {
                    "Subject": "fake websites",
                    "SubjectType": "domain-name",
                    "Relation": "include",
                    "Object": "hr-wipro[.]com and hr-suncor[.]com",
                    "ObjectType": "domain-name",
                    "tactic": "Reconnaissance",
                    "technique": [
                        "T1598-Phishing for Information"
                    ],
                    "id": 24
                }
            ]
        },
        "Resource Development": {
            "T1587-Develop Capabilities": [
                {
                    "Subject": "The DNSpionage campaign",
                    "SubjectType": "campaign",
                    "Relation": "demonstrate",
                    "Object": "resource development",
                    "ObjectType": "course-of-action",
                    "tactic": "Resource Development",
                    "technique": [
                        "T1587-Develop Capabilities"
                    ],
                    "id": 25
                },
                {
                    "Subject": "attackers",
                    "SubjectType": "threat-actor",
                    "Relation": "create",
                    "Object": "two fake websites",
                    "ObjectType": "infrastructure",
                    "tactic": "Resource Development",
                    "technique": [
                        "T1587-Develop Capabilities"
                    ],
                    "id": 26
                },
                {
                    "Subject": "attackers",
                    "SubjectType": "threat-actor",
                    "Relation": "generate",
                    "Object": "Let's Encrypt certificates",
                    "ObjectType": "X.509 Certificate",
                    "tactic": "Resource Development",
                    "technique": [
                        "T1587-Develop Capabilities"
                    ],
                    "id": 27
                },
                {
                    "Subject": "redirected domains",
                    "SubjectType": "domain-name",
                    "Relation": "use",
                    "Object": "Let's Encrypt certificates",
                    "ObjectType": "X.509 Certificate",
                    "tactic": "Resource Development",
                    "technique": [
                        "T1587-Develop Capabilities"
                    ],
                    "id": 28
                }
            ]
        },
        "Impact": {
            "T1565-Data Manipulation": [
                {
                    "Subject": "attackers",
                    "SubjectType": "threat-actor",
                    "Relation": "intercept",
                    "Object": "traffic",
                    "ObjectType": "network traffic",
                    "tactic": "Impact",
                    "technique": [
                        "T1565-Data Manipulation"
                    ],
                    "id": 30
                },
                {
                    "Subject": "attackers",
                    "SubjectType": "threat-actor",
                    "Relation": "harvest",
                    "Object": "email and VPN credentials",
                    "ObjectType": "artifacts",
                    "tactic": "Impact",
                    "technique": [
                        "T1565-Data Manipulation"
                    ],
                    "id": 31
                }
            ],
            "T1561-Disk Wipe": [
                {
                    "Subject": "attackers",
                    "SubjectType": "threat-actor",
                    "Relation": "enable",
                    "Object": "further attacks or blackmail",
                    "ObjectType": "attack-pattern",
                    "tactic": "Impact",
                    "technique": [
                        "T1561-Disk Wipe"
                    ],
                    "id": 32
                }
            ]
        }
    }
}