# question_generation/grouped_attackseqs/389.json
{
    "file_name": "The Enigmatic Roma225 Campaign",
    "tactic_label": true,
    "rewrite": {
        "Initial Access": "The Cybaze-Yoroi ZLab researchers discovered a malware campaign targeting the Italian automotive sector, delivered via a phishing email impersonating a senior partner of the Brazilian law firm 'Veirano Advogados'. The email contained a PowerPoint add-in document with a '.ppa' extension, which included an auto-open VBA macro designed to execute further malicious activities.",
        "Execution": "The VBA macro within the '.ppa' file used the 'mshta.exe' tool to download and execute a dropper from a Blogspot URL. The downloaded web page, while appearing benign, concealed VBScript code in an invisible blog post. This script was falsely attributed to 'Microsoft Corp.' and was responsible for executing a series of malicious actions, including storing a base64 encoded version of the 'RevengeRAT' payload in the registry, decoding and executing the payload, and creating and executing another VBScript to download a new payload from 'hxxp://cdtmaster.com[.]br'.",
        "Persistence": "The malware established persistence by creating a scheduled task to run the 'mshta.exe' utility with the URL 'hxxps://pocasideiascdt.blogspot[.]com/' every two hours, ensuring continuous execution of the malicious code.",
        "Command and Control": "The 'RevengeRAT' payload, once executed, established communication with its command and control servers, sending victim machine information to two different C2 destinations: 'office365update[.]duckdns.org' and 'systen32.ddns[.]net'. If one server was down, the malware would fall back to the other. Both C2 servers were unresponsive at the time of writing, but their IP resolutions suggested that the attacker's infrastructure could be located in Canada and Brazil.",
        "Others": "The 'Document.exe' file, characterized by a Pokemon Megaball image as its icon, was hosted at 'cdtmaster.com[.]br' and downloaded by the 'Z3j.vbs' script. Its sole purpose was to deploy and run the 'Outlook.exe' payload. The 'Outlook.exe' payload, when executed, did not show any immediate signs of activity but opened a listening TCP socket on localhost. The researchers were still analyzing the 'Outlook.exe' sample to determine its true behavior. The attack could not be attributed to a specific threat actor, but there were common TTPs with previous reports, such as the use of shared infrastructure as a drop-server and the use of popular RATs as backdoors. The domain 'cdtmaster.com[.]br' also hosted other suspicious files, including a variant of 'njRAT'. The indicators of compromise included various URLs and IP addresses associated with the campaign, as well as a registry key used for persistence. Yara rules were provided for detection of related malware samples."
    },
    "technique_label": true,
    "triplet_groups": {
        "Initial Access": {
            "T1566-Phishing": [
                {
                    "Subject": "Cybaze-Yoroi ZLab researchers",
                    "SubjectType": "author",
                    "Relation": "discovered",
                    "Object": "malware campaign",
                    "ObjectType": "campaign",
                    "tactic": "Initial Access",
                    "technique": [
                        "T1566-Phishing"
                    ],
                    "id": 1
                },
                {
                    "Subject": "malware campaign",
                    "SubjectType": "campaign",
                    "Relation": "target",
                    "Object": "Italian automotive sector",
                    "ObjectType": "identity",
                    "tactic": "Initial Access",
                    "technique": [
                        "T1566-Phishing"
                    ],
                    "id": 2
                },
                {
                    "Subject": "phishing email",
                    "SubjectType": "Email Message",
                    "Relation": "impersonate",
                    "Object": "senior partner of Veirano Advogados",
                    "ObjectType": "identity",
                    "tactic": "Initial Access",
                    "technique": [
                        "T1566-Phishing"
                    ],
                    "id": 3
                },
                {
                    "Subject": "phishing email",
                    "SubjectType": "Email Message",
                    "Relation": "deliver",
                    "Object": "PowerPoint add-in document",
                    "ObjectType": "file",
                    "tactic": "Initial Access",
                    "technique": [
                        "T1566-Phishing"
                    ],
                    "id": 4
                },
                {
                    "Subject": "PowerPoint add-in document",
                    "SubjectType": "file",
                    "Relation": "have",
                    "Object": "auto-open VBA macro",
                    "ObjectType": "artifact",
                    "tactic": "Initial Access",
                    "technique": [
                        "T1566-Phishing"
                    ],
                    "id": 5
                }
            ]
        },
        "Execution": {
            "T1059-Command and Scripting Interpreter": [
                {
                    "Subject": "The VBA macro",
                    "SubjectType": "artifact",
                    "Relation": "use",
                    "Object": "mshta.exe",
                    "ObjectType": "tool",
                    "tactic": "Execution",
                    "technique": [
                        "T1059-Command and Scripting Interpreter"
                    ],
                    "id": 6
                },
                {
                    "Subject": "mshta.exe",
                    "SubjectType": "tool",
                    "Relation": "download",
                    "Object": "dropper",
                    "ObjectType": "artifact",
                    "tactic": "Execution",
                    "technique": [
                        "T1059-Command and Scripting Interpreter"
                    ],
                    "id": 7
                },
                {
                    "Subject": "VBScript code",
                    "SubjectType": "artifact",
                    "Relation": "store",
                    "Object": "RevengeRAT payload",
                    "ObjectType": "Malware",
                    "tactic": "Execution",
                    "technique": [
                        "T1059-Command and Scripting Interpreter"
                    ],
                    "id": 8
                },
                {
                    "Subject": "VBScript code",
                    "SubjectType": "artifact",
                    "Relation": "execute",
                    "Object": "RevengeRAT payload",
                    "ObjectType": "Malware",
                    "tactic": "Execution",
                    "technique": [
                        "T1059-Command and Scripting Interpreter"
                    ],
                    "id": 9
                },
                {
                    "Subject": "VBScript code",
                    "SubjectType": "artifact",
                    "Relation": "download",
                    "Object": "new payload",
                    "ObjectType": "artifact",
                    "tactic": "Execution",
                    "technique": [
                        "T1059-Command and Scripting Interpreter"
                    ],
                    "id": 10
                }
            ]
        },
        "Persistence": {
            "T1053-Scheduled Task/Job": [
                {
                    "Subject": "The malware",
                    "SubjectType": "Malware",
                    "Relation": "establish persistence by creating",
                    "Object": "scheduled task",
                    "ObjectType": "course-of-action",
                    "tactic": "Persistence",
                    "technique": [
                        "T1053-Scheduled Task/Job"
                    ],
                    "id": 11
                },
                {
                    "Subject": "scheduled task",
                    "SubjectType": "course-of-action",
                    "Relation": "run",
                    "Object": "mshta.exe",
                    "ObjectType": "tool",
                    "tactic": "Persistence",
                    "technique": [
                        "T1053-Scheduled Task/Job"
                    ],
                    "id": 12
                }
            ]
        },
        "Command and Control": {
            "T1071-Application Layer Protocol": [
                {
                    "Subject": "RevengeRAT",
                    "SubjectType": "Malware",
                    "Relation": "establish communication with",
                    "Object": "command and control servers",
                    "ObjectType": "infrastructure",
                    "tactic": "Command and Control",
                    "technique": [
                        "T1071-Application Layer Protocol"
                    ],
                    "id": 14
                },
                {
                    "Subject": "RevengeRAT",
                    "SubjectType": "Malware",
                    "Relation": "send information to",
                    "Object": "office365update[.]duckdns.org",
                    "ObjectType": "domain-name",
                    "tactic": "Command and Control",
                    "technique": [
                        "T1071-Application Layer Protocol"
                    ],
                    "id": 15
                },
                {
                    "Subject": "RevengeRAT",
                    "SubjectType": "Malware",
                    "Relation": "send information to",
                    "Object": "systen32.ddns[.]net",
                    "ObjectType": "domain-name",
                    "tactic": "Command and Control",
                    "technique": [
                        "T1071-Application Layer Protocol"
                    ],
                    "id": 16
                },
                {
                    "Subject": "attacker's infrastructure",
                    "SubjectType": "infrastructure",
                    "Relation": "located at",
                    "Object": "Canada",
                    "ObjectType": "Location",
                    "tactic": "Command and Control",
                    "technique": [
                        "T1071-Application Layer Protocol"
                    ],
                    "id": 17
                },
                {
                    "Subject": "attacker's infrastructure",
                    "SubjectType": "infrastructure",
                    "Relation": "located at",
                    "Object": "Brazil",
                    "ObjectType": "Location",
                    "tactic": "Command and Control",
                    "technique": [
                        "T1071-Application Layer Protocol"
                    ],
                    "id": 18
                }
            ]
        }
    }
}