# question_generation/grouped_attackseqs/4.json
{
    "file_name": "Cisco Talos Intelligence Group - Comprehensive Threat Intelligence_ Mustang Panda deploys a new wave of malware targeting Europe",
    "tactic_label": true,
    "rewrite": {
        "Initial Access": "In February 2022, Cisco Talos observed Mustang Panda, a China-based threat actor, initiating phishing campaigns against European entities, including Russian organizations. The phishing emails contained malicious lures disguised as official European Union reports on the conflict in Ukraine and its effects on NATO countries, as well as fake Ukrainian government reports, which upon interaction, downloaded malware onto the victims' machines.",
        "Execution": "Mustang Panda employed various methods to execute their malware. They used benign executables to side-load malicious DLLs, which acted as loaders for the PlugX implant. The group also deployed custom stagers, meterpreter-based shellcode downloaders, and reverse shells. These methods were used to establish long-term persistence on infected endpoints for espionage purposes.",
        "Persistence": "Mustang Panda's tactics for maintaining persistence included the use of benign executables that loaded malicious DLLs, which in turn activated the PlugX RAT implant. They also used stagers and reverse shells to ensure continued access to compromised systems. The group utilized scheduled tasks and registry modifications to maintain their foothold.",
        "Command and Control": "Mustang Panda's command and control techniques involved the use of downloaders that fetched various infection artifacts, leading to the deployment of the PlugX implant. They also used stagers that downloaded additional shellcode from a C2 IP address, which was then executed on the infected endpoint.",
        "Others": "Mustang Panda, also known as RedDelta or Bronze President, has been active since at least 2012, targeting entities worldwide, including government organizations, think tanks, NGOs, and religious institutions. The group has extensively targeted Asian countries, such as the Taiwanese government, Hong Kong activists, NGOs in Mongolia and Tibet, and telecommunication firms in Myanmar and India. They have evolved their tactics and implants over the years to conduct espionage and information theft, using summit- and conference-themed lures in Asia and Europe. The group has also continuously evolved its delivery mechanisms, including maldocs, shortcut files, malicious archives, and downloaders. Cisco Talos Intelligence Group provides various solutions to detect and block threats from Mustang Panda."
    },
    "technique_label": true,
    "triplet_groups": {
        "Initial Access": {
            "T1566-Phishing": [
                {
                    "Subject": "Mustang Panda",
                    "SubjectType": "threat-actor",
                    "Relation": "initiate",
                    "Object": "phishing campaigns",
                    "ObjectType": "campaign",
                    "tactic": "Initial Access",
                    "technique": [
                        "T1566-Phishing"
                    ],
                    "id": 1
                },
                {
                    "Subject": "phishing emails",
                    "SubjectType": "Email Message",
                    "Relation": "contain",
                    "Object": "malicious lures",
                    "ObjectType": "artifacts",
                    "tactic": "Initial Access",
                    "technique": [
                        "T1566-Phishing"
                    ],
                    "id": 2
                },
                {
                    "Subject": "malicious lures",
                    "SubjectType": "artifacts",
                    "Relation": "disguised-as",
                    "Object": "official European Union reports",
                    "ObjectType": "artifacts",
                    "tactic": "Initial Access",
                    "technique": [
                        "T1566-Phishing"
                    ],
                    "id": 3
                },
                {
                    "Subject": "malicious lures",
                    "SubjectType": "artifacts",
                    "Relation": "disguised-as",
                    "Object": "fake Ukrainian government reports",
                    "ObjectType": "artifacts",
                    "tactic": "Initial Access",
                    "technique": [
                        "T1566-Phishing"
                    ],
                    "id": 4
                },
                {
                    "Subject": "malicious lures",
                    "SubjectType": "artifacts",
                    "Relation": "download",
                    "Object": "malware",
                    "ObjectType": "Malware",
                    "tactic": "Initial Access",
                    "technique": [
                        "T1566-Phishing"
                    ],
                    "id": 5
                }
            ]
        },
        "Execution": {
            "T1129-Shared Modules": [
                {
                    "Subject": "Mustang Panda",
                    "SubjectType": "threat-actor",
                    "Relation": "side-load",
                    "Object": "malicious DLLs",
                    "ObjectType": "file",
                    "tactic": "Execution",
                    "technique": [
                        "T1129-Shared Modules"
                    ],
                    "id": 7
                },
                {
                    "Subject": "malicious DLLs",
                    "SubjectType": "file",
                    "Relation": "act-as",
                    "Object": "loaders",
                    "ObjectType": "tool",
                    "tactic": "Execution",
                    "technique": [
                        "T1129-Shared Modules"
                    ],
                    "id": 8
                },
                {
                    "Subject": "loaders",
                    "SubjectType": "tool",
                    "Relation": "load",
                    "Object": "PlugX implant",
                    "ObjectType": "Malware",
                    "tactic": "Execution",
                    "technique": [
                        "T1129-Shared Modules"
                    ],
                    "id": 9
                }
            ]
        },
        "Persistence": {
            "T1053-Scheduled Task/Job": [
                {
                    "Subject": "Mustang Panda",
                    "SubjectType": "threat-actor",
                    "Relation": "utilize",
                    "Object": "scheduled tasks",
                    "ObjectType": "technique",
                    "tactic": "Persistence",
                    "technique": [
                        "T1053-Scheduled Task/Job"
                    ],
                    "id": 18
                }
            ],
            "T1547-Boot or Logon Autostart Execution": [
                {
                    "Subject": "Mustang Panda",
                    "SubjectType": "threat-actor",
                    "Relation": "utilize",
                    "Object": "registry modifications",
                    "ObjectType": "technique",
                    "tactic": "Persistence",
                    "technique": [
                        "T1547-Boot or Logon Autostart Execution"
                    ],
                    "id": 19
                }
            ]
        },
        "Command and Control": {
            "T1105-Ingress Tool Transfer": [
                {
                    "Subject": "Mustang Panda",
                    "SubjectType": "threat-actor",
                    "Relation": "use",
                    "Object": "downloaders",
                    "ObjectType": "tool",
                    "tactic": "Command and Control",
                    "technique": [
                        "T1105-Ingress Tool Transfer"
                    ],
                    "id": 20
                },
                {
                    "Subject": "downloaders",
                    "SubjectType": "tool",
                    "Relation": "fetch",
                    "Object": "infection artifacts",
                    "ObjectType": "artifacts",
                    "tactic": "Command and Control",
                    "technique": [
                        "T1105-Ingress Tool Transfer"
                    ],
                    "id": 21
                },
                {
                    "Subject": "stagers",
                    "SubjectType": "tool",
                    "Relation": "download",
                    "Object": "additional shellcode",
                    "ObjectType": "artifacts",
                    "tactic": "Command and Control",
                    "technique": [
                        "T1105-Ingress Tool Transfer"
                    ],
                    "id": 23
                },
                {
                    "Subject": "stagers",
                    "SubjectType": "tool",
                    "Relation": "execute",
                    "Object": "shellcode",
                    "ObjectType": "artifacts",
                    "tactic": "Command and Control",
                    "technique": [
                        "T1105-Ingress Tool Transfer"
                    ],
                    "id": 24
                }
            ]
        }
    }
}