# cti_preprocessing/filtered_attackseqs/134.json
{
    "triplets": [
        {
            "Subject": "targeted campaign",
            "SubjectType": "campaign",
            "Relation": "carried-out-against",
            "Object": "Israeli organizations",
            "ObjectType": "identity",
            "tactic": "Initial Access",
            "technique": [
                "T1190-Exploit Public-Facing Application"
            ],
            "id": 1
        },
        {
            "Subject": "Microsoft Word documents",
            "SubjectType": "file",
            "Relation": "exploit",
            "Object": "CVE-2017-0199",
            "ObjectType": "vulnerability",
            "tactic": "Initial Access",
            "technique": [
                "T1190-Exploit Public-Facing Application"
            ],
            "id": 2
        },
        {
            "Subject": "Microsoft Word documents",
            "SubjectType": "file",
            "Relation": "install",
            "Object": "Helminth Trojan",
            "ObjectType": "Malware",
            "tactic": "Initial Access",
            "technique": [
                "Others"
            ],
            "id": 3
        },
        {
            "Subject": "The Helminth Trojan agent",
            "SubjectType": "Malware",
            "Relation": "use",
            "Object": "various evasive mechanisms",
            "ObjectType": "Technique",
            "tactic": "Defense Evasion",
            "technique": [
                "T1027-Obfuscated Files or Information"
            ],
            "id": 4
        },
        {
            "Subject": "The attack",
            "SubjectType": "Campaign",
            "Relation": "utilize",
            "Object": "a sophisticated version of the Mimikatz script",
            "ObjectType": "Tool",
            "tactic": "Defense Evasion",
            "technique": [
                "T1027-Obfuscated Files or Information"
            ],
            "id": 5
        },
        {
            "Subject": "bat script",
            "SubjectType": "tool",
            "Relation": "variant-of",
            "Object": "Mimikatz",
            "ObjectType": "Malware",
            "tactic": "Credential Access",
            "technique": [
                "T1555-Credentials from Password Stores"
            ],
            "id": 6
        },
        {
            "Subject": "bat script",
            "SubjectType": "tool",
            "Relation": "use",
            "Object": "scheduled task",
            "ObjectType": "artifact",
            "tactic": "Credential Access",
            "technique": [
                "Others"
            ],
            "id": 7
        },
        {
            "Subject": "DnE1.Ps1 file",
            "SubjectType": "file",
            "Relation": "communicate-with",
            "Object": "C&C server",
            "ObjectType": "domain-name",
            "tactic": "Command and Control",
            "technique": [
                "T1071-Application Layer Protocol"
            ],
            "id": 8
        },
        {
            "Subject": "DnE1.Ps1 file",
            "SubjectType": "file",
            "Relation": "download",
            "Object": "bat script",
            "ObjectType": "file",
            "tactic": "Command and Control",
            "technique": [
                "T1105-Ingress Tool Transfer"
            ],
            "id": 9
        },
        {
            "Subject": "The attack",
            "SubjectType": "campaign",
            "Relation": "utilize",
            "Object": "vulnerability exploit",
            "ObjectType": "technique",
            "tactic": "Execution",
            "technique": [
                "T1059-Command and Scripting Interpreter"
            ],
            "id": 10
        },
        {
            "Subject": "The .hta file",
            "SubjectType": "file",
            "Relation": "deliver",
            "Object": "mshta.exe",
            "ObjectType": "process",
            "tactic": "Execution",
            "technique": [
                "T1059-Command and Scripting Interpreter"
            ],
            "id": 11
        },
        {
            "Subject": "mshta.exe",
            "SubjectType": "process",
            "Relation": "execute",
            "Object": "The .hta file",
            "ObjectType": "file",
            "tactic": "Execution",
            "technique": [
                "T1059-Command and Scripting Interpreter"
            ],
            "id": 12
        },
        {
            "Subject": "The attack",
            "SubjectType": "campaign",
            "Relation": "disable",
            "Object": "warning message",
            "ObjectType": "artifact",
            "tactic": "Execution",
            "technique": [
                "T1059-Command and Scripting Interpreter"
            ],
            "id": 13
        },
        {
            "Subject": "Helminth Trojan agent",
            "SubjectType": "Malware",
            "Relation": "achieve persistence by adding",
            "Object": "GoogleUpdateTasksMachineUI",
            "ObjectType": "process",
            "tactic": "Persistence",
            "technique": [
                "T1053-Scheduled Task/Job"
            ],
            "id": 14
        },
        {
            "Subject": "Helminth Trojan agent",
            "SubjectType": "Malware",
            "Relation": "execute",
            "Object": "VBS file",
            "ObjectType": "file",
            "tactic": "Persistence",
            "technique": [
                "T1053-Scheduled Task/Job"
            ],
            "id": 15
        }
    ],
    "file_name": "Iranian Fileless Attack Infiltrates Israeli Organizations",
    "tactic_label": true,
    "rewrite": {
        "Initial Access": "From April 19-24, 2017, a targeted campaign was carried out against Israeli organizations. The attack was delivered through compromised email accounts at Ben-Gurion University, using Microsoft Word documents that exploited the CVE-2017-0199 vulnerability to install a fileless variant of the Helminth Trojan agent. The Word document was sent as a reply to legitimate emails from those accounts and was propagated to more than 250 individuals in different Israeli companies.",
        "Execution": "The attack utilized a vulnerability exploit rather than macros. The .hta file delivered by the embedded link in the Word document was executed by mshta.exe, which typically shows a warning message. However, the attack disabled this message by sending an 'Enter' command to the warning window.",
        "Persistence": "The Helminth Trojan agent achieved persistence by adding a scheduled task named 'GoogleUpdateTasksMachineUI', which executes a VBS file every 3 minutes. The files are installed in the Public\\Libraries\\RecordedTV\\ folder.",
        "Defense Evasion": "The Helminth Trojan agent used various evasive mechanisms, including the regeneration of its signatures on the infected computer to slow down detection and remediation. The attack also utilized a sophisticated version of the Mimikatz script to gather information while avoiding detection.",
        "Credential Access": "The bat script used in the attack is a customized version of Mimikatz, designed to gather information from the computer and the network. It includes commands to handle non-ASCII characters and validate the scheduled task added by the persistence mechanism.",
        "Command and Control": "The DnE1.Ps1 file communicates with a C&C server using DNS exchange queries, which is difficult to block. The server, vpsupdate[.]tk, was first registered on April 16, 2017. The script's goal is to download a bat script, execute it, and upload the results back to the C&C, clearing traces afterward.",
        "Others": "The C&C servers are still active and will be listed in the report. The attack is attributed to an Iranian hacker group responsible for the OilRig malware campaigns. The group is speculated to be sponsored by the Iranian government. The attack's sophistication has increased, with improvements in infrastructure, delivery servers, and the Helminth agent. The attackers set up the attack quickly, speculating that the window of opportunity between patch release and patch rollout was still open."
    },
    "technique_label": true
}