# cti_preprocessing/filtered_attackseqs/303.json
{
    "triplets": [
        {
            "Subject": "T9000 malware",
            "SubjectType": "malware",
            "Relation": "variant-of",
            "Object": "T5000 family",
            "ObjectType": "malware",
            "tactic": "Initial Access",
            "technique": [
                "T1190-Exploit Public-Facing Application"
            ],
            "id": 1
        },
        {
            "Subject": "T9000 malware",
            "SubjectType": "malware",
            "Relation": "delivered-via",
            "Object": "RTF file",
            "ObjectType": "file",
            "tactic": "Initial Access",
            "technique": [
                "T1190-Exploit Public-Facing Application"
            ],
            "id": 2
        },
        {
            "Subject": "RTF file",
            "SubjectType": "file",
            "Relation": "exploit",
            "Object": "CVE-2012-1856",
            "ObjectType": "vulnerability",
            "tactic": "Initial Access",
            "technique": [
                "T1190-Exploit Public-Facing Application"
            ],
            "id": 3
        },
        {
            "Subject": "RTF file",
            "SubjectType": "file",
            "Relation": "exploit",
            "Object": "CVE-2015-1641",
            "ObjectType": "vulnerability",
            "tactic": "Initial Access",
            "technique": [
                "T1190-Exploit Public-Facing Application"
            ],
            "id": 4
        },
        {
            "Subject": "initial shellcode stage",
            "SubjectType": "tool",
            "Relation": "responsible-for",
            "Object": "locating and executing secondary shellcode stub",
            "ObjectType": "tool",
            "tactic": "Initial Access",
            "technique": [
                "T1190-Exploit Public-Facing Application"
            ],
            "id": 5
        },
        {
            "Subject": "secondary shellcode stub",
            "SubjectType": "tool",
            "Relation": "decrypt and execute",
            "Object": "embedded payload",
            "ObjectType": "file",
            "tactic": "Initial Access",
            "technique": [
                "T1190-Exploit Public-Facing Application"
            ],
            "id": 6
        },
        {
            "Subject": "embedded payload",
            "SubjectType": "file",
            "Relation": "create",
            "Object": "temporary file",
            "ObjectType": "file",
            "tactic": "Initial Access",
            "technique": [
                "T1190-Exploit Public-Facing Application"
            ],
            "id": 7
        },
        {
            "Subject": "temporary file",
            "SubjectType": "file",
            "Relation": "execute",
            "Object": "WinExec",
            "ObjectType": "tool",
            "tactic": "Initial Access",
            "technique": [
                "T1190-Exploit Public-Facing Application"
            ],
            "id": 8
        },
        {
            "Subject": "malware",
            "SubjectType": "Malware",
            "Relation": "communicate-with",
            "Object": "C2 server",
            "ObjectType": "infrastrcucture",
            "tactic": "Command and Control",
            "technique": [
                "T1071-Application Layer Protocol"
            ],
            "id": 9
        },
        {
            "Subject": "malware",
            "SubjectType": "Malware",
            "Relation": "send",
            "Object": "encrypted traffic",
            "ObjectType": "network traffic",
            "tactic": "Command and Control",
            "technique": [
                "T1573-Encrypted Channel"
            ],
            "id": 10
        },
        {
            "Subject": "malware",
            "SubjectType": "Malware",
            "Relation": "support",
            "Object": "proxy usage",
            "ObjectType": "feature",
            "tactic": "Command and Control",
            "technique": [
                "T1090-Proxy"
            ],
            "id": 11
        },
        {
            "Subject": "T9000",
            "SubjectType": "Malware",
            "Relation": "involved",
            "Object": "igfxtray.exe",
            "ObjectType": "file",
            "tactic": "Execution",
            "technique": [
                "T1129-Shared Modules"
            ],
            "id": 12
        },
        {
            "Subject": "T9000",
            "SubjectType": "Malware",
            "Relation": "sideloaded",
            "Object": "hccutils.dll",
            "ObjectType": "file",
            "tactic": "Execution",
            "technique": [
                "T1129-Shared Modules"
            ],
            "id": 13
        },
        {
            "Subject": "T9000",
            "SubjectType": "Malware",
            "Relation": "performed",
            "Object": "security product checks",
            "ObjectType": "Technique",
            "tactic": "Execution",
            "technique": [
                "Others"
            ],
            "id": 14
        },
        {
            "Subject": "T9000",
            "SubjectType": "Malware",
            "Relation": "used",
            "Object": "techniques for persistence",
            "ObjectType": "Technique",
            "tactic": "Execution",
            "technique": [
                "Others"
            ],
            "id": 15
        },
        {
            "Subject": "T9000",
            "SubjectType": "Malware",
            "Relation": "designed to",
            "Object": "evade detection",
            "ObjectType": "Tactic",
            "tactic": "Defense Evasion",
            "technique": [
                "T1497-Virtualization/Sandbox Evasion"
            ],
            "id": 16
        },
        {
            "Subject": "T9000",
            "SubjectType": "Malware",
            "Relation": "use",
            "Object": "multi-stage installation process",
            "ObjectType": "Technique",
            "tactic": "Defense Evasion",
            "technique": [
                "Others"
            ],
            "id": 17
        },
        {
            "Subject": "T9000",
            "SubjectType": "Malware",
            "Relation": "create",
            "Object": "mutexes and named events",
            "ObjectType": "artifacts",
            "tactic": "Defense Evasion",
            "technique": [
                "T1564-Hide Artifacts"
            ],
            "id": 18
        },
        {
            "Subject": "T9000",
            "SubjectType": "Malware",
            "Relation": "modify",
            "Object": "HKLM\\Software\\Microsoft\\Windows\\CurrentVersion\\Run",
            "ObjectType": "registry",
            "tactic": "Persistence",
            "technique": [
                "T1547-Boot or Logon Autostart Execution"
            ],
            "id": 19
        },
        {
            "Subject": "T9000",
            "SubjectType": "Malware",
            "Relation": "modify",
            "Object": "HKLM\\Software\\Microsoft\\Windows NT\\CurrentVersion\\Windows",
            "ObjectType": "registry",
            "tactic": "Persistence",
            "technique": [
                "T1547-Boot or Logon Autostart Execution"
            ],
            "id": 20
        },
        {
            "Subject": "T9000",
            "SubjectType": "Malware",
            "Relation": "inject",
            "Object": "ResN32.dll",
            "ObjectType": "file",
            "tactic": "Persistence",
            "technique": [
                "T1547-Boot or Logon Autostart Execution"
            ],
            "id": 21
        },
        {
            "Subject": "ResN32.dll",
            "SubjectType": "file",
            "Relation": "load",
            "Object": "explorer.exe",
            "ObjectType": "process",
            "tactic": "Persistence",
            "technique": [
                "T1547-Boot or Logon Autostart Execution"
            ],
            "id": 22
        },
        {
            "Subject": "T9000",
            "SubjectType": "Malware",
            "Relation": "capture",
            "Object": "data",
            "ObjectType": "artifact",
            "tactic": "Collection",
            "technique": [
                "T1119-Automated Collection",
                "T1025-Data from Removable Media"
            ],
            "id": 23
        },
        {
            "Subject": "T9000",
            "SubjectType": "Malware",
            "Relation": "steal",
            "Object": "files",
            "ObjectType": "artifact",
            "tactic": "Collection",
            "technique": [
                "T1119-Automated Collection",
                "T1025-Data from Removable Media"
            ],
            "id": 24
        },
        {
            "Subject": "T9000",
            "SubjectType": "Malware",
            "Relation": "target",
            "Object": "Skype users",
            "ObjectType": "identity",
            "tactic": "Collection",
            "technique": [
                "T1119-Automated Collection",
                "T1123-Audio Capture",
                "T1125-Video Capture",
                "T1113-Screen Capture"
            ],
            "id": 25
        },
        {
            "Subject": "T9000",
            "SubjectType": "Malware",
            "Relation": "exfiltrate-to",
            "Object": "C2 server",
            "ObjectType": "infrastructure",
            "tactic": "Exfiltration",
            "technique": [
                "T1041-Exfiltration Over C2 Channel"
            ],
            "id": 26
        }
    ],
    "file_name": "PaloAlto_t9000-advanced-modular-backdoor-uses-complex-anti-analysis-techniques.2016.02.04",
    "tactic_label": true,
    "rewrite": {
        "Initial Access": "The T9000 malware, a variant of the T5000 family, was initially delivered via an RTF file exploiting CVE-2012-1856 and CVE-2015-1641. Upon activation, a multi-stage execution process began, starting with an initial shellcode stage responsible for locating and executing a secondary shellcode stub. This stub then decrypted and executed an embedded payload, creating a temporary file in the %TEMP% folder and executing it using WinExec.",
        "Execution": "The execution of T9000 involved a legitimate Microsoft executable, igfxtray.exe, which sideloaded a malicious DLL, hccutils.dll. This DLL performed security product checks and initiated the next stage of execution depending on the victim's system properties and installed security products. The malware used three techniques for persistence, including setting registry keys and injecting the ResN32.dll library into the explorer.exe process.",
        "Persistence": "T9000 established persistence through various methods, including modifying registry keys such as HKLM\\Software\\Microsoft\\Windows\\CurrentVersion\\Run and HKLM\\Software\\Microsoft\\Windows NT\\CurrentVersion\\Windows, to ensure the malicious DLLs would run after reboot or load with every user mode process spawned. Additionally, it injected the ResN32.dll library into the explorer.exe process for persistent execution.",
        "Defense Evasion": "T9000 was designed to evade detection by customizing its installation to avoid 24 identified security products. It used a multi-stage installation process with checks at each stage to detect analysis by security researchers. The malware also created mutexes and named events to ensure only one instance was running and to manage its components.",
        "Collection": "T9000 was configured to automatically capture data about the infected system and steal files of specific types stored on removable media. It also targeted Skype users, capturing encrypted data, screenshots of specific applications, and Skype communications including video calls, audio calls, and chat messages.",
        "Command and Control": "The malware communicated with a C2 server using a binary protocol over a specified port, sending encrypted traffic that included a beacon and collected information from the victim machine. The C2 server's IP address and port were extracted from the configuration file, and the malware supported proxy usage for victims behind web proxies.",
        "Exfiltration": "T9000 exfiltrated various types of data, including screenshots, Skype communications, and files from removable drives. The data was compressed, encrypted, and stored in specific directories or sent to the C2 server. The malware also logged file and directory operations, clipboard data, and encrypted network traffic to collect sensitive information.",
        "Others": "The T9000 malware was observed in targeted attacks against organizations in the United States, with functionality indicating broad user targeting. It was part of the T5000 malware family and known for its extensive evasion techniques. Palo Alto Networks' Threat Prevention, Traps, WildFire, and AutoFocus provided protection against T9000/T5000 attacks. The malware's author took significant measures to avoid detection and analysis by the security community."
    },
    "technique_label": true
}