# cti_preprocessing/filtered_attackseqs/320.json
{
    "triplets": [
        {
            "Subject": "Shuckworm",
            "SubjectType": "threat-actor",
            "Relation": "initiate",
            "Object": "attack chain",
            "ObjectType": "campaign",
            "tactic": "Initial Access",
            "technique": [
                "T1566-Phishing"
            ],
            "id": 1
        },
        {
            "Subject": "Shuckworm",
            "SubjectType": "threat-actor",
            "Relation": "use",
            "Object": "Pterodo",
            "ObjectType": "Malware",
            "tactic": "Initial Access",
            "technique": [
                "Others"
            ],
            "id": 2
        },
        {
            "Subject": "Shuckworm",
            "SubjectType": "threat-actor",
            "Relation": "create",
            "Object": "scheduled tasks",
            "ObjectType": "artifacts",
            "tactic": "Persistence",
            "technique": [
                "T1547-Boot or Logon Autostart Execution",
                "T1053-Scheduled Task/Job"
            ],
            "id": 3
        },
        {
            "Subject": "Shuckworm",
            "SubjectType": "threat-actor",
            "Relation": "use",
            "Object": "mshta.exe",
            "ObjectType": "file",
            "tactic": "Defense Evasion",
            "technique": [
                "T1218-System Binary Proxy Execution"
            ],
            "id": 4
        },
        {
            "Subject": "Shuckworm",
            "SubjectType": "threat-actor",
            "Relation": "use",
            "Object": "scheduled tasks",
            "ObjectType": "Technique",
            "tactic": "Defense Evasion",
            "technique": [
                "T1053-Scheduled Task"
            ],
            "id": 5
        },
        {
            "Subject": "Shuckworm",
            "SubjectType": "threat-actor",
            "Relation": "use",
            "Object": "VBS scripts",
            "ObjectType": "file",
            "tactic": "Defense Evasion",
            "technique": [
                "Others"
            ],
            "id": 6
        },
        {
            "Subject": "Shuckworm",
            "SubjectType": "threat-actor",
            "Relation": "deploy",
            "Object": "VBS backdoor",
            "ObjectType": "Malware",
            "tactic": "Collection",
            "technique": [
                "Others"
            ],
            "id": 7
        },
        {
            "Subject": "Shuckworm",
            "SubjectType": "threat-actor",
            "Relation": "execute",
            "Object": "VBS scripts",
            "ObjectType": "artifacts",
            "tactic": "Collection",
            "technique": [
                "Others"
            ],
            "id": 8
        },
        {
            "Subject": "Shuckworm",
            "SubjectType": "threat-actor",
            "Relation": "execute",
            "Object": "VBS scripts",
            "ObjectType": "file",
            "tactic": "Execution",
            "technique": [
                "T1059-Command and Scripting Interpreter"
            ],
            "id": 9
        },
        {
            "Subject": "Shuckworm",
            "SubjectType": "threat-actor",
            "Relation": "execute",
            "Object": "HTA files",
            "ObjectType": "file",
            "tactic": "Execution",
            "technique": [
                "T1059-Command and Scripting Interpreter"
            ],
            "id": 10
        },
        {
            "Subject": "Shuckworm",
            "SubjectType": "threat-actor",
            "Relation": "create",
            "Object": "scheduled tasks",
            "ObjectType": "Technique",
            "tactic": "Execution",
            "technique": [
                "T1053-Scheduled Task/Job"
            ],
            "id": 11
        },
        {
            "Subject": "Shuckworm",
            "SubjectType": "threat-actor",
            "Relation": "use",
            "Object": "VNC client",
            "ObjectType": "tool",
            "tactic": "Execution",
            "technique": [
                "Others"
            ],
            "id": 12
        },
        {
            "Subject": "Shuckworm",
            "SubjectType": "threat-actor",
            "Relation": "establish",
            "Object": "command and control",
            "ObjectType": "network traffic",
            "tactic": "Command and Control",
            "technique": [
                "T1071-Application Layer Protocol"
            ],
            "id": 13
        },
        {
            "Subject": "Shuckworm",
            "SubjectType": "threat-actor",
            "Relation": "download",
            "Object": "additional payloads",
            "ObjectType": "file",
            "tactic": "Command and Control",
            "technique": [
                "T1071-Application Layer Protocol",
                "T1105-Ingress Tool Transfer"
            ],
            "id": 14
        },
        {
            "Subject": "Shuckworm",
            "SubjectType": "threat-actor",
            "Relation": "use",
            "Object": "C&C servers",
            "ObjectType": "infrastructure",
            "tactic": "Command and Control",
            "technique": [
                "T1071-Application Layer Protocol"
            ],
            "id": 15
        }
    ],
    "file_name": "Shuckworm_APT",
    "tactic_label": true,
    "rewrite": {
        "Initial Access": "The attack chain initiated by Shuckworm against a Ukrainian organization began on July 14, 2021, with the opening of a malicious document, likely delivered via phishing email. The document contained a custom backdoor known as Pterodo, which was executed using a suspicious command.",
        "Execution": "Following the initial compromise, Shuckworm executed various VBS scripts and HTA files to maintain their presence and execute their backdoor. They also created scheduled tasks to ensure persistence and executed commands via scripts. On July 16, additional variants of Pterodo were executed, and on July 28, a VNC client was installed as the ultimate payload for the attack.",
        "Persistence": "Shuckworm established persistence by creating scheduled tasks that executed malicious VBS files every few minutes. These tasks ensured that their backdoor would remain active even after system reboots.",
        "Defense Evasion": "Shuckworm utilized mshta.exe to bypass application control solutions and browser security settings. They also used scheduled tasks and VBS scripts to evade detection and maintain their foothold.",
        "Collection": "Between July 29 and August 18, Shuckworm deployed multiple variants of their custom VBS backdoor and executed VBS scripts. It is unclear if the documents accessed during this period were part of the attackers' efforts to collect and exfiltrate sensitive information.",
        "Command and Control": "Shuckworm established command and control by communicating with their servers using HTTP requests and downloading additional payloads. They used various C&C servers, including deep-pitched.enarto.ru, arianat.ru, and mucoris.ru, to control the compromised systems and download tools.",
        "Others": "Symantec's Threat Hunter Team, part of Broadcom Software, has been tracking Shuckworm's cyber-espionage activities against Ukraine. Shuckworm, also known as Gamaredon or Armageddon, has been active since at least 2013 and has recently increased the sophistication of their attacks. The Security Service of Ukraine (SSU) reported that Shuckworm has been using living-off-the-land tools to steal credentials and move laterally on victim networks. Symantec's investigation revealed seven files used by Shuckworm in their attacks, all of which are 7-zip SFX self-extracting binaries. These files included various VBS files, a VNC client, and a Process Explorer binary. The IOCs from the investigation showed patterns such as C&C IPs belonging to a short list of hosting providers and suspected malicious files found in specific directories with names starting with the letter 'd'. The SSU report from November 2021 attributed over 5,000 attacks to Shuckworm since 2014, targeting more than 1,500 Ukrainian government systems."
    },
    "technique_label": true
}