# cti_preprocessing/filtered_attackseqs/400.json
{
    "triplets": [
        {
            "Subject": "Malicious executable",
            "SubjectType": "file",
            "Relation": "execute",
            "Object": "malicious HTA file",
            "ObjectType": "file",
            "tactic": "Execution",
            "technique": [
                "T1059-Command and Scripting Interpreter",
                "T1204-User Execution"
            ],
            "id": 1
        },
        {
            "Subject": "HTA script",
            "SubjectType": "file",
            "Relation": "fetch",
            "Object": "VBS code",
            "ObjectType": "artifact",
            "tactic": "Execution",
            "technique": [
                "T1059-Command and Scripting Interpreter"
            ],
            "id": 2
        },
        {
            "Subject": "VBS code",
            "SubjectType": "artifact",
            "Relation": "download",
            "Object": "final payload",
            "ObjectType": "Malware",
            "tactic": "Execution",
            "technique": [
                "T1059-Command and Scripting Interpreter"
            ],
            "id": 3
        },
        {
            "Subject": "Final payload",
            "SubjectType": "Malware",
            "Relation": "variant-of",
            "Object": "banking Trojan",
            "ObjectType": "Malware",
            "tactic": "Execution",
            "technique": [
                "Others"
            ],
            "id": 4
        },
        {
            "Subject": "Persistence",
            "SubjectType": "Tactic",
            "Relation": "achieve",
            "Object": "registry keys",
            "ObjectType": "registry",
            "tactic": "Persistence",
            "technique": [
                "T1547-Boot or Logon Autostart Execution"
            ],
            "id": 5
        },
        {
            "Subject": "registry keys",
            "SubjectType": "registry",
            "Relation": "create",
            "Object": "Software\\Microsoft\\Windows\\CurrentVersion\\Run",
            "ObjectType": "directory",
            "tactic": "Persistence",
            "technique": [
                "T1547-Boot or Logon Autostart Execution"
            ],
            "id": 6
        },
        {
            "Subject": "VBS script",
            "SubjectType": "file",
            "Relation": "set to persist",
            "Object": "infection chain",
            "ObjectType": "artifacts",
            "tactic": "Persistence",
            "technique": [
                "T1547-Boot or Logon Autostart Execution"
            ],
            "id": 7
        },
        {
            "Subject": "VBS script",
            "SubjectType": "file",
            "Relation": "download",
            "Object": "ZIP archive",
            "ObjectType": "file",
            "tactic": "Persistence",
            "technique": [
                "T1547-Boot or Logon Autostart Execution"
            ],
            "id": 8
        },
        {
            "Subject": "campaigns",
            "SubjectType": "Campaign",
            "Relation": "use",
            "Object": "public cloud infrastructure",
            "ObjectType": "infrastructure",
            "tactic": "Defense Evasion",
            "technique": [
                "T1578-Modify Cloud Compute Infrastructure"
            ],
            "id": 9
        },
        {
            "Subject": "campaigns",
            "SubjectType": "Campaign",
            "Relation": "abuse",
            "Object": "legitimate signed binaries",
            "ObjectType": "software",
            "tactic": "Defense Evasion",
            "technique": [
                "T1218-System Binary Proxy Execution"
            ],
            "id": 10
        },
        {
            "Subject": "The malware",
            "SubjectType": "Malware",
            "Relation": "communicate-with",
            "Object": "C2 server",
            "ObjectType": "infrastructure",
            "tactic": "Command and Control",
            "technique": [
                "T1571-Non-Standard Port"
            ],
            "id": 11
        },
        {
            "Subject": "Metamorfo campaigns",
            "SubjectType": "campaign",
            "Relation": "begin-with",
            "Object": "phishing emails",
            "ObjectType": "attack-pattern",
            "tactic": "Initial Access",
            "technique": [
                "T1566-Phishing"
            ],
            "id": 12
        },
        {
            "Subject": "Metamorfo campaigns",
            "SubjectType": "campaign",
            "Relation": "redirect",
            "Object": "victims to cloud storage sites like GitHub, Dropbox, or Google Drive",
            "ObjectType": "infrastructure",
            "tactic": "Initial Access",
            "technique": [
                "Others"
            ],
            "id": 13
        },
        {
            "Subject": "Metamorfo campaigns",
            "SubjectType": "campaign",
            "Relation": "deliver",
            "Object": "ZIP file containing a malicious executable",
            "ObjectType": "file",
            "tactic": "Initial Access",
            "technique": [
                "Others"
            ],
            "id": 14
        },
        {
            "Subject": "banking Trojan",
            "SubjectType": "Malware",
            "Relation": "search for",
            "Object": "antivirus software",
            "ObjectType": "software",
            "tactic": "Discovery",
            "technique": [
                "T1518-Software Discovery"
            ],
            "id": 15
        },
        {
            "Subject": "banking Trojan",
            "SubjectType": "Malware",
            "Relation": "search for",
            "Object": "banking software",
            "ObjectType": "software",
            "tactic": "Discovery",
            "technique": [
                "T1518-Software Discovery"
            ],
            "id": 16
        },
        {
            "Subject": "banking Trojan",
            "SubjectType": "Malware",
            "Relation": "collect information about",
            "Object": "hostname",
            "ObjectType": "artifact",
            "tactic": "Discovery",
            "technique": [
                "T1082-System Information Discovery",
                "T1016-System Network Configuration Discovery"
            ],
            "id": 17
        },
        {
            "Subject": "banking Trojan",
            "SubjectType": "Malware",
            "Relation": "collect information about",
            "Object": "operating system version",
            "ObjectType": "artifact",
            "tactic": "Discovery",
            "technique": [
                "T1082-System Information Discovery",
                "T1016-System Network Configuration Discovery"
            ],
            "id": 18
        },
        {
            "Subject": "banking Trojan",
            "SubjectType": "Malware",
            "Relation": "send to",
            "Object": "hardcoded domain",
            "ObjectType": "domain-name",
            "tactic": "Discovery",
            "technique": [
                "Others"
            ],
            "id": 19
        },
        {
            "Subject": "The malware",
            "SubjectType": "Malware",
            "Relation": "capture",
            "Object": "screenshots and mouse clicks",
            "ObjectType": "artifacts",
            "tactic": "Collection",
            "technique": [
                "T1113-Screen Capture"
            ],
            "id": 20
        },
        {
            "Subject": "screenshots",
            "SubjectType": "artifacts",
            "Relation": "save-as",
            "Object": ".jpg images",
            "ObjectType": "file",
            "tactic": "Collection",
            "technique": [
                "T1113-Screen Capture"
            ],
            "id": 21
        },
        {
            "Subject": "screenshots",
            "SubjectType": "artifacts",
            "Relation": "locate-at",
            "Object": "a specific directory",
            "ObjectType": "directory",
            "tactic": "Collection",
            "technique": [
                "T1113-Screen Capture"
            ],
            "id": 22
        }
    ],
    "file_name": "metamorfo-campaign-targeting-brazilian-users_html",
    "tactic_label": true,
    "rewrite": {
        "Initial Access": "The Metamorfo campaigns began with phishing emails containing HTML attachments that redirected victims to cloud storage sites like GitHub, Dropbox, or Google Drive to download a ZIP file containing a malicious executable. This executable, when run, would continue the infection chain.",
        "Execution": "The malicious executable from the initial access stage was an HTA script compiled into an executable that, when run, would execute a malicious HTA file extracted to %temp% and run by mshta.exe. This HTA script contained VBS code that fetched additional VBS code encoded in base64 from a remote server, which then downloaded the final payload, a banking Trojan, from another remote server.",
        "Persistence": "Persistence was achieved by creating registry keys in Software\\Microsoft\\Windows\\CurrentVersion\\Run using a string extracted from a file named car.dat. Additionally, a VBS script capable of recreating the infection chain was set to persist and download the same ZIP archive.",
        "Defense Evasion": "The campaigns utilized various techniques to evade detection, including the use of public cloud infrastructure to deliver stages of the malicious payload and the abuse of legitimate signed binaries to load malicious code.",
        "Discovery": "The banking Trojan would search for folders related to antivirus software and banking software on the victim's system. It would also collect information about the hostname and operating system version, which was sent to a hardcoded domain using a hardcoded User-Agent string.",
        "Collection": "The malware captured screenshots and mouse clicks while the user was browsing banking sites. These screenshots were saved as .jpg images in a specific directory.",
        "Command and Control": "The C2 server was selected based on a string in a file named 'id'. The malware communicated with the C2 server over raw TCP on port 9999, following a specific command pattern. The C2 server's infrastructure revealed the number of potential victims for the campaign.",
        "Others": "The Metamorfo campaigns targeted Brazilian users with banking Trojans delivered through malspam. The campaigns were notable for their multi-stage infection chains and the use of cloud storage sites to host malicious payloads. The attackers also used legitimate Windows tools for DLL side-loading and employed various evasion techniques. The campaigns were analyzed by FireEye Labs, which provided insights into the tactics and techniques used by the attackers."
    },
    "technique_label": true
}