# question_generation/filtered_attackseqs/13.json
{
    "triplets": [
        {
            "Subject": "RTF document",
            "SubjectType": "file",
            "Relation": "exploit",
            "Object": "CVE-2018-0798",
            "ObjectType": "vulnerability",
            "tactic": "Initial Access",
            "technique": [
                "T1190-Exploit Public-Facing Application"
            ],
            "id": 1
        },
        {
            "Subject": "OSEA54d.tmp",
            "SubjectType": "file",
            "Relation": "located-at",
            "Object": "C:\\Users\\[username]\\AppData\\Local\\Temp directory",
            "ObjectType": "directory",
            "tactic": "Initial Access",
            "technique": [
                "T1190-Exploit Public-Facing Application"
            ],
            "id": 2
        },
        {
            "Subject": "Google-Updates",
            "SubjectType": "artifact",
            "Relation": "use",
            "Object": "msiexec",
            "ObjectType": "software",
            "tactic": "Persistence",
            "technique": [
                "T1053-Scheduled Task/Job"
            ],
            "id": 3
        },
        {
            "Subject": "Google-Updates",
            "SubjectType": "artifact",
            "Relation": "execute",
            "Object": "[ComputerName].tmp",
            "ObjectType": "file",
            "tactic": "Persistence",
            "technique": [
                "T1053-Scheduled Task/Job"
            ],
            "id": 4
        },
        {
            "Subject": "Google-Update",
            "SubjectType": "artifact",
            "Relation": "run",
            "Object": "Photo.vbs",
            "ObjectType": "file",
            "tactic": "Persistence",
            "technique": [
                "T1053-Scheduled Task/Job"
            ],
            "id": 5
        },
        {
            "Subject": "attackers",
            "SubjectType": "threat-actor",
            "Relation": "use",
            "Object": "CertUtil.exe",
            "ObjectType": "tool",
            "tactic": "Defense Evasion",
            "technique": [
                "T1036-Masquerading"
            ],
            "id": 6
        },
        {
            "Subject": "CertUtil.exe",
            "SubjectType": "tool",
            "Relation": "disguise-as",
            "Object": "GoogleUpdate.exe",
            "ObjectType": "file",
            "tactic": "Defense Evasion",
            "technique": [
                "T1036-Masquerading"
            ],
            "id": 7
        },
        {
            "Subject": "CertUtil.exe",
            "SubjectType": "tool",
            "Relation": "download",
            "Object": "additional payloads",
            "ObjectType": "artifact",
            "tactic": "Defense Evasion",
            "technique": [
                "T1070-Indicator Removal"
            ],
            "id": 8
        },
        {
            "Subject": "The command and control mechanism",
            "SubjectType": "tool",
            "Relation": "use",
            "Object": "an IP address (167.71.237.100)",
            "ObjectType": "ipv4-addr",
            "tactic": "Command and Control",
            "technique": [
                "T1571-Non-Standard Port"
            ],
            "id": 9
        },
        {
            "Subject": "OSEA54d.tmp file",
            "SubjectType": "file",
            "Relation": "drop",
            "Object": "GoogleUpdate.exe",
            "ObjectType": "file",
            "tactic": "Execution",
            "technique": [
                "T1059-Command and Scripting Interpreter"
            ],
            "id": 10
        },
        {
            "Subject": "GoogleUpdate.exe",
            "SubjectType": "file",
            "Relation": "impersonate",
            "Object": "CertUtil.exe",
            "ObjectType": "software",
            "tactic": "Execution",
            "technique": [
                "T1059-Command and Scripting Interpreter"
            ],
            "id": 11
        },
        {
            "Subject": "Photo.vbs",
            "SubjectType": "file",
            "Relation": "facilitate",
            "Object": "download",
            "ObjectType": "attack-pattern",
            "tactic": "Execution",
            "technique": [
                "T1059-Command and Scripting Interpreter"
            ],
            "id": 12
        }
    ],
    "file_name": "RANCOR APT_ Suspected targeted attacks against South East Asia",
    "tactic_label": true,
    "rewrite": {
        "Initial Access": "On August 28, 2019, an RTF document was created and uploaded to VirusTotal on September 11, 2019. This document, written in Khmer, appeared to be sent by the Cabinet of Cambodia and listed names of government officials. It was designed to exploit CVE-2018-0798, a stack buffer overflow vulnerability, to execute a file named 'OSEA54d.tmp' located in the 'C:\\Users\\[username]\\AppData\\Local\\Temp' directory.",
        "Execution": "The 'OSEA54d.tmp' file subsequently dropped 'GoogleUpdate.exe' (hash: 7b973145f7e1b59330ca4dd1f86b3d55) into 'C:\\Windows\\System32\\spool\\drivers\\color\\'. This executable was actually CertUtil.exe, a legitimate Microsoft utility. Additionally, a .vbs script named 'Photo.vbs' was created in the same directory to facilitate the download of a file named after the user's computer name with a .png extension, which was then saved as a .tmp file in the user's local temp folder.",
        "Persistence": "Two scheduled tasks were created for persistence: 'Google-Updates' and 'Google-Update'. The former used 'msiexec' to execute the downloaded payload '[ComputerName].tmp' every minute, while the latter ran 'Photo.vbs' every two minutes.",
        "Defense Evasion": "The attackers used a legitimate Microsoft utility, CertUtil.exe, disguised as 'GoogleUpdate.exe' to evade detection and download additional payloads.",
        "Command and Control": "The command and control mechanism involved the use of an IP address (167.71.237.100) to serve a second stage payload, which was not retrieved during the analysis due to HTTP 404 responses.",
        "Others": "The RTF document targeted South East Asian entities and was suspected to be part of a campaign by Rancor APT, a Chinese APT group identified by PaloAlto in 2017. The document contained a list of names and appeared official, likely as a lure for espionage purposes. The final payload was not available for analysis, but it was likely to be DDKONG or PLAINTEE malware, which serve as backdoors for file manipulation and command execution. The TTPs observed suggest a high degree of confidence that the activity is related to Rancor APT."
    },
    "technique_label": true
}