# question_generation/filtered_attackseqs/218.json
{
    "triplets": [
        {
            "Subject": "Dexphot",
            "SubjectType": "Malware",
            "Relation": "ensure",
            "Object": "persistence",
            "ObjectType": "Tactic",
            "tactic": "Persistence",
            "technique": [
                "T1053-Scheduled Task/Job"
            ],
            "id": 1
        },
        {
            "Subject": "Dexphot",
            "SubjectType": "Malware",
            "Relation": "use",
            "Object": "schtasks.exe",
            "ObjectType": "tool",
            "tactic": "Persistence",
            "technique": [
                "T1053-Scheduled Task/Job"
            ],
            "id": 2
        },
        {
            "Subject": "Dexphot",
            "SubjectType": "Malware",
            "Relation": "employ",
            "Object": "Scheduled Task",
            "ObjectType": "Technique",
            "tactic": "Persistence",
            "technique": [
                "T1053-Scheduled Task/Job"
            ],
            "id": 3
        },
        {
            "Subject": "Dexphot",
            "SubjectType": "Malware",
            "Relation": "employ",
            "Object": "Signed Binary Proxy Execution",
            "ObjectType": "Technique",
            "tactic": "Persistence",
            "technique": [
                "Others"
            ],
            "id": 4
        },
        {
            "Subject": "tasks",
            "SubjectType": "artifact",
            "Relation": "call",
            "Object": "msiexec.exe",
            "ObjectType": "tool",
            "tactic": "Persistence",
            "technique": [
                "T1053-Scheduled Task/Job"
            ],
            "id": 5
        },
        {
            "Subject": "Dexphot",
            "SubjectType": "Malware",
            "Relation": "use",
            "Object": "SoftwareBundler:Win32/ICLoader",
            "ObjectType": "tool",
            "tactic": "Defense Evasion",
            "technique": [
                "T1027-Obfuscated Files or Information"
            ],
            "id": 6
        },
        {
            "Subject": "Dexphot",
            "SubjectType": "Malware",
            "Relation": "use",
            "Object": "living-off-the-land techniques",
            "ObjectType": "Tactic",
            "tactic": "Defense Evasion",
            "technique": [
                "Others"
            ],
            "id": 7
        },
        {
            "Subject": "Dexphot",
            "SubjectType": "Malware",
            "Relation": "use",
            "Object": "fileless techniques",
            "ObjectType": "Tactic",
            "tactic": "Defense Evasion",
            "technique": [
                "T1564-Hide Artifacts"
            ],
            "id": 8
        },
        {
            "Subject": "Dexphot",
            "SubjectType": "Malware",
            "Relation": "use",
            "Object": "process hollowing",
            "ObjectType": "Technique",
            "tactic": "Defense Evasion",
            "technique": [
                "T1055-Process Injection"
            ],
            "id": 9
        },
        {
            "Subject": "Dexphot attack",
            "SubjectType": "campaign",
            "Relation": "use",
            "Object": "installer",
            "ObjectType": "file",
            "tactic": "Execution",
            "technique": [
                "T1059-Command and Scripting Interpreter"
            ],
            "id": 10
        },
        {
            "Subject": "Dexphot attack",
            "SubjectType": "campaign",
            "Relation": "use",
            "Object": "MSI package file",
            "ObjectType": "file",
            "tactic": "Execution",
            "technique": [
                "T1059-Command and Scripting Interpreter"
            ],
            "id": 11
        },
        {
            "Subject": "Dexphot attack",
            "SubjectType": "campaign",
            "Relation": "use",
            "Object": "password-protected ZIP archive",
            "ObjectType": "file",
            "tactic": "Execution",
            "technique": [
                "T1059-Command and Scripting Interpreter"
            ],
            "id": 12
        },
        {
            "Subject": "Dexphot attack",
            "SubjectType": "campaign",
            "Relation": "use",
            "Object": "loader DLL",
            "ObjectType": "file",
            "tactic": "Execution",
            "technique": [
                "T1059-Command and Scripting Interpreter"
            ],
            "id": 13
        },
        {
            "Subject": "Dexphot attack",
            "SubjectType": "campaign",
            "Relation": "use",
            "Object": "encrypted data file",
            "ObjectType": "file",
            "tactic": "Execution",
            "technique": [
                "T1059-Command and Scripting Interpreter"
            ],
            "id": 14
        },
        {
            "Subject": "Dexphot attack",
            "SubjectType": "campaign",
            "Relation": "target",
            "Object": "msiexec.exe",
            "ObjectType": "process",
            "tactic": "Execution",
            "technique": [
                "T1059-Command and Scripting Interpreter"
            ],
            "id": 15
        },
        {
            "Subject": "Dexphot attack",
            "SubjectType": "campaign",
            "Relation": "target",
            "Object": "unzip.exe",
            "ObjectType": "process",
            "tactic": "Execution",
            "technique": [
                "T1059-Command and Scripting Interpreter"
            ],
            "id": 16
        },
        {
            "Subject": "Dexphot attack",
            "SubjectType": "campaign",
            "Relation": "target",
            "Object": "rundll32.exe",
            "ObjectType": "process",
            "tactic": "Execution",
            "technique": [
                "T1059-Command and Scripting Interpreter"
            ],
            "id": 17
        },
        {
            "Subject": "Dexphot attack",
            "SubjectType": "campaign",
            "Relation": "target",
            "Object": "schtasks.exe",
            "ObjectType": "process",
            "tactic": "Execution",
            "technique": [
                "T1059-Command and Scripting Interpreter"
            ],
            "id": 18
        },
        {
            "Subject": "Dexphot attack",
            "SubjectType": "campaign",
            "Relation": "target",
            "Object": "powershell.exe",
            "ObjectType": "process",
            "tactic": "Execution",
            "technique": [
                "T1059-Command and Scripting Interpreter"
            ],
            "id": 19
        },
        {
            "Subject": "Dexphot attack",
            "SubjectType": "campaign",
            "Relation": "target",
            "Object": "svchost.exe",
            "ObjectType": "process",
            "tactic": "Execution",
            "technique": [
                "T1059-Command and Scripting Interpreter"
            ],
            "id": 20
        },
        {
            "Subject": "Dexphot attack",
            "SubjectType": "campaign",
            "Relation": "target",
            "Object": "tracert.exe",
            "ObjectType": "process",
            "tactic": "Execution",
            "technique": [
                "T1059-Command and Scripting Interpreter"
            ],
            "id": 21
        },
        {
            "Subject": "Dexphot attack",
            "SubjectType": "campaign",
            "Relation": "target",
            "Object": "setup.exe",
            "ObjectType": "process",
            "tactic": "Execution",
            "technique": [
                "T1059-Command and Scripting Interpreter"
            ],
            "id": 22
        }
    ],
    "file_name": "Insights from one year of tracking a polymorphic threat",
    "tactic_label": true,
    "rewrite": {
        "Execution": "The Dexphot attack, first detected in October 2018, involved a complex execution stage where it wrote five key files to disk: an installer with two URLs, an MSI package file, a password-protected ZIP archive, a loader DLL, and an encrypted data file containing three additional executables. These executables were loaded into system processes via process hollowing, targeting legitimate system processes such as msiexec.exe, unzip.exe, rundll32.exe, schtasks.exe, and powershell.exe, and later svchost.exe, tracert.exe, and setup.exe.",
        "Persistence": "Dexphot ensured persistence through dual monitoring services that checked the status of all three malicious processes and re-infected the device if any were terminated. It also used schtasks.exe to create scheduled tasks for updates and re-infection, employing techniques like Scheduled Task and Signed Binary Proxy Execution. The tasks called msiexec.exe as a proxy to run the malicious code and updated Dexphot's components upon system reboot and every 90 or 110 minutes while the system was running.",
        "Defense Evasion": "Dexphot employed multiple layers of security evasion, including the use of SoftwareBundler:Win32/ICLoader to drop and run the installer, living-off-the-land techniques, and obfuscated batch scripts to check for antivirus products. It also used fileless techniques and process hollowing to hide malicious code within legitimate system processes, making detection and remediation difficult.",
        "Others": "Microsoft Defender ATP Research Team, led by Hazel Kim, closely monitored the Dexphot threat for over a year, observing its evolution and the attackers' efforts to upgrade the malware, target new processes, and work around defensive measures. Microsoft Defender ATP's behavioral blocking and containment capabilities were crucial in stopping the Dexphot campaign. The team's insights into the goals and motivations of Dexphot's authors contributed to a broader understanding of cybercriminal behavior."
    },
    "technique_label": true
}