# question_generation/filtered_attackseqs/266.json
{
    "triplets": [
        {
            "Subject": "FireEye",
            "SubjectType": "CTI vender",
            "Relation": "identified",
            "Object": "phishing campaign",
            "ObjectType": "campaign",
            "tactic": "Initial Access",
            "technique": [
                "T1566-Phishing"
            ],
            "id": 1
        },
        {
            "Subject": "APT34",
            "SubjectType": "threat-actor",
            "Relation": "conducted",
            "Object": "phishing campaign",
            "ObjectType": "campaign",
            "tactic": "Initial Access",
            "technique": [
                "T1566-Phishing"
            ],
            "id": 2
        },
        {
            "Subject": "APT34",
            "SubjectType": "threat-actor",
            "Relation": "masquerade",
            "Object": "Cambridge University member",
            "ObjectType": "identity",
            "tactic": "Initial Access",
            "technique": [
                "T1566-Phishing"
            ],
            "id": 3
        },
        {
            "Subject": "APT34",
            "SubjectType": "threat-actor",
            "Relation": "use",
            "Object": "LinkedIn",
            "ObjectType": "software",
            "tactic": "Initial Access",
            "technique": [
                "T1566-Phishing"
            ],
            "id": 4
        },
        {
            "Subject": "Rebecca Watts",
            "SubjectType": "identity",
            "Relation": "deliver",
            "Object": "malicious documents",
            "ObjectType": "file",
            "tactic": "Initial Access",
            "technique": [
                "T1566-Phishing"
            ],
            "id": 5
        },
        {
            "Subject": "System.doc",
            "SubjectType": "file",
            "Relation": "identified-as",
            "Object": "TONEDEAF",
            "ObjectType": "Malware",
            "tactic": "Execution",
            "technique": [
                "T1059-Command and Scripting Interpreter"
            ],
            "id": 6
        },
        {
            "Subject": "TONEDEAF",
            "SubjectType": "Malware",
            "Relation": "communicate-with",
            "Object": "C2 server",
            "ObjectType": "infrastructure",
            "tactic": "Execution",
            "technique": [
                "T1059-Command and Scripting Interpreter"
            ],
            "id": 7
        },
        {
            "Subject": "ERFT-Details.xls",
            "SubjectType": "file",
            "Relation": "create",
            "Object": "windows update check",
            "ObjectType": "task",
            "tactic": "Persistence",
            "technique": [
                "T1053-Scheduled Task/Job",
                "T1137-Office Application Startup"
            ],
            "id": 8
        },
        {
            "Subject": "windows update check",
            "SubjectType": "task",
            "Relation": "run",
            "Object": "System Manager.exe",
            "ObjectType": "file",
            "tactic": "Persistence",
            "technique": [
                "T1053-Scheduled Task/Job",
                "T1137-Office Application Startup"
            ],
            "id": 9
        },
        {
            "Subject": "System.doc",
            "SubjectType": "file",
            "Relation": "rename to",
            "Object": "System Manager.exe",
            "ObjectType": "file",
            "tactic": "Persistence",
            "technique": [
                "T1137-Office Application Startup"
            ],
            "id": 10
        },
        {
            "Subject": "FireEye's ExploitGuard engine",
            "SubjectType": "tool",
            "Relation": "detected and stopped",
            "Object": "exploit attempt by Microsoft Excel",
            "ObjectType": "attack-pattern",
            "tactic": "Defense Evasion",
            "technique": [
                "Others"
            ],
            "id": 11
        },
        {
            "Subject": "APT34",
            "SubjectType": "threat-actor",
            "Relation": "add",
            "Object": "PICKPOCKET",
            "ObjectType": "Malware",
            "tactic": "Credential Access",
            "technique": [
                "T1555-Credentials from Password Stores"
            ],
            "id": 12
        },
        {
            "Subject": "APT34",
            "SubjectType": "threat-actor",
            "Relation": "use",
            "Object": "PICKPOCKET",
            "ObjectType": "Malware",
            "tactic": "Credential Access",
            "technique": [
                "T1555-Credentials from Password Stores"
            ],
            "id": 13
        },
        {
            "Subject": "APT34",
            "SubjectType": "threat-actor",
            "Relation": "use",
            "Object": "VALUEVAULT",
            "ObjectType": "Malware",
            "tactic": "Credential Access",
            "technique": [
                "T1555-Credentials from Password Stores"
            ],
            "id": 14
        },
        {
            "Subject": "APT34",
            "SubjectType": "threat-actor",
            "Relation": "use",
            "Object": "LONGWATCH",
            "ObjectType": "Malware",
            "tactic": "Credential Access",
            "technique": [
                "T1056-Input Capture"
            ],
            "id": 15
        },
        {
            "Subject": "TONEDEAF",
            "SubjectType": "Malware",
            "Relation": "communicate-with",
            "Object": "C2 server",
            "ObjectType": "infrastructure",
            "tactic": "Command and Control",
            "technique": [
                "T1071-Application Layer Protocol"
            ],
            "id": 16
        },
        {
            "Subject": "C2 domain",
            "SubjectType": "domain-name",
            "Relation": "identified-as",
            "Object": "offlineearthquake[.]com",
            "ObjectType": "domain-name",
            "tactic": "Command and Control",
            "technique": [
                "T1071-Application Layer Protocol"
            ],
            "id": 17
        },
        {
            "Subject": "Malware",
            "SubjectType": "Malware",
            "Relation": "retrieve",
            "Object": "sys_id",
            "ObjectType": "artifacts",
            "tactic": "Command and Control",
            "technique": [
                "T1071-Application Layer Protocol"
            ],
            "id": 18
        }
    ],
    "file_name": "Hard Pass_ Declining APT34\u2019s Invite to Join Their Professional Network",
    "tactic_label": true,
    "rewrite": {
        "Initial Access": "In late June 2019, FireEye identified a phishing campaign conducted by APT34, an Iranian-nexus threat actor. The campaign involved masquerading as a member of Cambridge University to gain victims' trust to open malicious documents and using LinkedIn to deliver these documents. The targeted employee received a file via a LinkedIn message from 'Rebecca Watts', who claimed to be 'Research Staff at University of Cambridge'. The conversation began with the solicitation of resumes for potential job opportunities.",
        "Execution": "The malicious file named System.doc, identified as a new malware family TONEDEAF by FireEye, was a Windows Portable Executable (PE) despite having a 'doc' file extension. TONEDEAF is a backdoor that communicates with a command and control (C2) server using HTTP GET and POST requests, supporting system information collection, file uploading and downloading, and arbitrary shell command execution.",
        "Persistence": "The spreadsheet ERFT-Details.xls created a scheduled task named 'windows update check' that runs the file C:\\Users\\<user_name>\\.templates\\System Manager.exe every minute. Upon closing the spreadsheet, a final VBA function renames System.doc to System Manager.exe.",
        "Defense Evasion": "FireEye's ExploitGuard engine, part of their Endpoint Security, detected and stopped the exploit attempt by Microsoft Excel. The ExploitGuard monitors application behavior for anomalies used by threat actors to subvert traditional detection mechanisms. The offending applications can be sandboxed or terminated, preventing an exploit from reaching its next programmed step.",
        "Credential Access": "APT34 added three new malware families to their arsenal, including PICKPOCKET, a malware exclusively observed in use by APT34 for browser credential theft. Additionally, VALUEVAULT, a Golang compiled version of the 'Windows Vault Password Dumper' was used for browser credential theft, and LONGWATCH, a keylogger that outputs keystrokes to a log.txt file in the Windows temp folder.",
        "Command and Control": "TONEDEAF communicated with a single C2 server using HTTP GET and POST requests. The C2 domain identified was offlineearthquake[.]com. The malware retrieves the system and current user names to create a 'sys_id' used in subsequent requests to track infected target activity.",
        "Others": "FireEye's Managed Defense Security Operations Center received an exploit detection alert on June 19, 2019. The Managed Defense SOC analyzed the alert and identified the malicious file System.doc. FireEye's Intelligence and Advanced Practices teams performed a wider search across global visibility, identifying additional artifacts and activity from APT34 actors at other victim organizations. FireEye's platform successfully thwarted the attempted intrusion, stopping the new malware variant. FireEye also initiated a Community Protection Event titled 'Geopolitical Spotlight: Iran' to update customers with new discoveries, activity, and detection efforts related to this campaign and other recent activity from Iranian-nexus threat actors."
    },
    "technique_label": true
}