# cti_preprocessing/grouped_attackseqs/252.json
{
    "file_name": "Tracking Turla_ New backdoor delivered via Armenian watering holes _ WeLiveSecurity",
    "tactic_label": true,
    "rewrite": {
        "Initial Access": "ESET researchers discovered a watering hole operation targeting Armenian websites, using a fake Adobe Flash update lure to deliver malware dubbed NetFlash and PyFlash. The operation is attributed to Turla, an espionage group targeting governmental and military organizations. Compromised websites include armconsul[.]ru, mnp.nkr[.]am, aiisa[.]am, and adgf[.]am, with some belonging to the Armenian government. The campaign involved inserting malicious JavaScript into these websites, such as appending obfuscated code to jquery-migrate.min.js on mnp.nkr[.]am. The malicious code would load external JavaScript from skategirlchina[.]com/wp-includes/data_from_db_top.php, which fingerprints the visitor's browser and potentially delivers a second-stage malicious JavaScript. The campaign ceased delivering malicious scripts by the end of November 2019.",
        "Execution": "The watering hole operation by Turla involved a social engineering trick with a fake Adobe Flash update warning to induce users to download and install malware. The fake installer, once executed, would install both a Turla malware variant and a legitimate Adobe Flash program, leading users to believe the update was legitimate.",
        "Persistence": "The NetFlash and PyFlash malware variants used by Turla establish persistence using Windows scheduled tasks. NetFlash downloads its second stage malware from a hardcoded URL and sets up a scheduled task for the new backdoor, while PyFlash persists by modifying the Windows task that regularly launches the malware.",
        "Discovery": "The PyFlash backdoor executes various OS-related commands such as systeminfo, tasklist, ipconfig, getmac, and arp to gather information about the machine and sends it to the C&C server.",
        "Command and Control": "PyFlash communicates with its hardcoded C&C server via HTTP, using AES-128 in CBC mode to encrypt all network communications. NetFlash uses port 80 for communication, while PyFlash uses port 8,000, and another NetFlash sample uses port 15,363.",
        "Exfiltration": "The output of PyFlash surveillance and C&C commands are exfiltrated using the C&C protocol.",
        "Others": "The Turla group has been active for over ten years, with previous campaigns including Mosquito and LightNeuron. The recent operation shares similarities with past Turla campaigns, particularly one uncovered in 2017. The JavaScript used in the current campaign is almost identical to that of the 2017 campaign, but with different targets and payloads. The operation relies on a well-known social engineering trick and has shifted payloads, likely to evade detection. The new payload, NetFlash, installs a backdoor called PyFlash, developed in Python. ESET will continue monitoring Turla activities and publish relevant information on their blog. Indicators of Compromise (IoCs) include compromised websites and C&C servers such as skategirlchina[.]com and various IP addresses. Samples of the malware with SHA-1 hashes and timestamps are also provided."
    },
    "technique_label": true,
    "triplet_groups": {
        "Initial Access": {
            "T1189-Drive-by Compromise": [
                {
                    "Subject": "ESET researchers",
                    "SubjectType": "CTI vender",
                    "Relation": "discovered",
                    "Object": "watering hole operation",
                    "ObjectType": "campaign",
                    "tactic": "Initial Access",
                    "technique": [
                        "T1189-Drive-by Compromise"
                    ],
                    "id": 1
                },
                {
                    "Subject": "Turla",
                    "SubjectType": "threat-actor",
                    "Relation": "attributed-to",
                    "Object": "watering hole operation",
                    "ObjectType": "campaign",
                    "tactic": "Initial Access",
                    "technique": [
                        "T1189-Drive-by Compromise"
                    ],
                    "id": 2
                },
                {
                    "Subject": "Turla",
                    "SubjectType": "threat-actor",
                    "Relation": "use",
                    "Object": "fake Adobe Flash update lure",
                    "ObjectType": "tool",
                    "tactic": "Initial Access",
                    "technique": [
                        "T1189-Drive-by Compromise"
                    ],
                    "id": 3
                },
                {
                    "Subject": "Turla",
                    "SubjectType": "threat-actor",
                    "Relation": "deliver",
                    "Object": "NetFlash",
                    "ObjectType": "Malware",
                    "tactic": "Initial Access",
                    "technique": [
                        "T1189-Drive-by Compromise"
                    ],
                    "id": 4
                },
                {
                    "Subject": "Turla",
                    "SubjectType": "threat-actor",
                    "Relation": "deliver",
                    "Object": "PyFlash",
                    "ObjectType": "Malware",
                    "tactic": "Initial Access",
                    "technique": [
                        "T1189-Drive-by Compromise"
                    ],
                    "id": 5
                },
                {
                    "Subject": "compromised websites",
                    "SubjectType": "infrastructure",
                    "Relation": "belong-to",
                    "Object": "Armenian government",
                    "ObjectType": "identity",
                    "tactic": "Initial Access",
                    "technique": [
                        "T1189-Drive-by Compromise"
                    ],
                    "id": 6
                },
                {
                    "Subject": "malicious JavaScript",
                    "SubjectType": "artifact",
                    "Relation": "insert-into",
                    "Object": "armconsul[.]ru",
                    "ObjectType": "domain-name",
                    "tactic": "Initial Access",
                    "technique": [
                        "T1189-Drive-by Compromise"
                    ],
                    "id": 7
                },
                {
                    "Subject": "malicious JavaScript",
                    "SubjectType": "artifact",
                    "Relation": "insert-into",
                    "Object": "mnp.nkr[.]am",
                    "ObjectType": "domain-name",
                    "tactic": "Initial Access",
                    "technique": [
                        "T1189-Drive-by Compromise"
                    ],
                    "id": 8
                },
                {
                    "Subject": "malicious JavaScript",
                    "SubjectType": "artifact",
                    "Relation": "insert-into",
                    "Object": "aiisa[.]am",
                    "ObjectType": "domain-name",
                    "tactic": "Initial Access",
                    "technique": [
                        "T1189-Drive-by Compromise"
                    ],
                    "id": 9
                },
                {
                    "Subject": "malicious JavaScript",
                    "SubjectType": "artifact",
                    "Relation": "insert-into",
                    "Object": "adgf[.]am",
                    "ObjectType": "domain-name",
                    "tactic": "Initial Access",
                    "technique": [
                        "T1189-Drive-by Compromise"
                    ],
                    "id": 10
                },
                {
                    "Subject": "malicious code",
                    "SubjectType": "artifact",
                    "Relation": "load-from",
                    "Object": "skategirlchina[.]com/wp-includes/data_from_db_top.php",
                    "ObjectType": "url",
                    "tactic": "Initial Access",
                    "technique": [
                        "T1189-Drive-by Compromise"
                    ],
                    "id": 11
                }
            ]
        },
        "Persistence": {
            "T1053-Scheduled Task/Job": [
                {
                    "Subject": "NetFlash",
                    "SubjectType": "Malware",
                    "Relation": "establish",
                    "Object": "Windows scheduled tasks",
                    "ObjectType": "artifact",
                    "tactic": "Persistence",
                    "technique": [
                        "T1053-Scheduled Task/Job"
                    ],
                    "id": 12
                },
                {
                    "Subject": "NetFlash",
                    "SubjectType": "Malware",
                    "Relation": "download",
                    "Object": "second stage malware",
                    "ObjectType": "Malware",
                    "tactic": "Persistence",
                    "technique": [
                        "T1053-Scheduled Task/Job"
                    ],
                    "id": 13
                },
                {
                    "Subject": "NetFlash",
                    "SubjectType": "Malware",
                    "Relation": "set up",
                    "Object": "scheduled task",
                    "ObjectType": "artifact",
                    "tactic": "Persistence",
                    "technique": [
                        "T1053-Scheduled Task/Job"
                    ],
                    "id": 14
                },
                {
                    "Subject": "PyFlash",
                    "SubjectType": "Malware",
                    "Relation": "persist",
                    "Object": "Windows task",
                    "ObjectType": "artifact",
                    "tactic": "Persistence",
                    "technique": [
                        "T1053-Scheduled Task/Job"
                    ],
                    "id": 15
                }
            ]
        },
        "Command and Control": {
            "T1071-Application Layer Protocol": [
                {
                    "Subject": "PyFlash",
                    "SubjectType": "Malware",
                    "Relation": "communicate-with",
                    "Object": "its hardcoded C&C server",
                    "ObjectType": "infrastrcucture",
                    "tactic": "Command and Control",
                    "technique": [
                        "T1071-Application Layer Protocol",
                        "T1573-Encrypted Channel"
                    ],
                    "id": 16
                },
                {
                    "Subject": "NetFlash",
                    "SubjectType": "Malware",
                    "Relation": "use",
                    "Object": "port 80",
                    "ObjectType": "network traffic",
                    "tactic": "Command and Control",
                    "technique": [
                        "T1071-Application Layer Protocol",
                        "T1571-Non-Standard Port"
                    ],
                    "id": 17
                },
                {
                    "Subject": "PyFlash",
                    "SubjectType": "Malware",
                    "Relation": "use",
                    "Object": "port 8,000",
                    "ObjectType": "network traffic",
                    "tactic": "Command and Control",
                    "technique": [
                        "T1071-Application Layer Protocol",
                        "T1571-Non-Standard Port"
                    ],
                    "id": 18
                },
                {
                    "Subject": "another NetFlash sample",
                    "SubjectType": "Malware",
                    "Relation": "use",
                    "Object": "port 15,363",
                    "ObjectType": "network traffic",
                    "tactic": "Command and Control",
                    "technique": [
                        "T1071-Application Layer Protocol",
                        "T1571-Non-Standard Port"
                    ],
                    "id": 19
                }
            ],
            "T1573-Encrypted Channel": [
                {
                    "Subject": "PyFlash",
                    "SubjectType": "Malware",
                    "Relation": "communicate-with",
                    "Object": "its hardcoded C&C server",
                    "ObjectType": "infrastrcucture",
                    "tactic": "Command and Control",
                    "technique": [
                        "T1071-Application Layer Protocol",
                        "T1573-Encrypted Channel"
                    ],
                    "id": 16
                }
            ],
            "T1571-Non-Standard Port": [
                {
                    "Subject": "NetFlash",
                    "SubjectType": "Malware",
                    "Relation": "use",
                    "Object": "port 80",
                    "ObjectType": "network traffic",
                    "tactic": "Command and Control",
                    "technique": [
                        "T1071-Application Layer Protocol",
                        "T1571-Non-Standard Port"
                    ],
                    "id": 17
                },
                {
                    "Subject": "PyFlash",
                    "SubjectType": "Malware",
                    "Relation": "use",
                    "Object": "port 8,000",
                    "ObjectType": "network traffic",
                    "tactic": "Command and Control",
                    "technique": [
                        "T1071-Application Layer Protocol",
                        "T1571-Non-Standard Port"
                    ],
                    "id": 18
                },
                {
                    "Subject": "another NetFlash sample",
                    "SubjectType": "Malware",
                    "Relation": "use",
                    "Object": "port 15,363",
                    "ObjectType": "network traffic",
                    "tactic": "Command and Control",
                    "technique": [
                        "T1071-Application Layer Protocol",
                        "T1571-Non-Standard Port"
                    ],
                    "id": 19
                }
            ]
        },
        "Execution": {
            "T1204-User Execution": [
                {
                    "Subject": "watering hole operation",
                    "SubjectType": "campaign",
                    "Relation": "use",
                    "Object": "social engineering trick",
                    "ObjectType": "technique",
                    "tactic": "Execution",
                    "technique": [
                        "T1204-User Execution"
                    ],
                    "id": 20
                },
                {
                    "Subject": "fake installer",
                    "SubjectType": "artifact",
                    "Relation": "install",
                    "Object": "Turla malware variant",
                    "ObjectType": "malware",
                    "tactic": "Execution",
                    "technique": [
                        "T1204-User Execution"
                    ],
                    "id": 21
                },
                {
                    "Subject": "fake installer",
                    "SubjectType": "artifact",
                    "Relation": "install",
                    "Object": "Adobe Flash program",
                    "ObjectType": "software",
                    "tactic": "Execution",
                    "technique": [
                        "T1204-User Execution"
                    ],
                    "id": 22
                }
            ]
        },
        "Discovery": {
            "T1082-System Information Discovery": [
                {
                    "Subject": "The PyFlash backdoor",
                    "SubjectType": "Malware",
                    "Relation": "execute",
                    "Object": "systeminfo, tasklist, ipconfig, getmac, arp",
                    "ObjectType": "command",
                    "tactic": "Discovery",
                    "technique": [
                        "T1082-System Information Discovery"
                    ],
                    "id": 23
                },
                {
                    "Subject": "The PyFlash backdoor",
                    "SubjectType": "Malware",
                    "Relation": "send",
                    "Object": "information",
                    "ObjectType": "data",
                    "tactic": "Discovery",
                    "technique": [
                        "T1082-System Information Discovery"
                    ],
                    "id": 24
                },
                {
                    "Subject": "information",
                    "SubjectType": "data",
                    "Relation": "be",
                    "Object": "to the C&C server",
                    "ObjectType": "infrastructure",
                    "tactic": "Discovery",
                    "technique": [
                        "T1082-System Information Discovery"
                    ],
                    "id": 25
                }
            ]
        },
        "Exfiltration": {
            "T1041-Exfiltration Over C2 Channel": [
                {
                    "Subject": "The output of PyFlash surveillance",
                    "SubjectType": "artifacts",
                    "Relation": "exfiltrate-to",
                    "Object": "C&C protocol",
                    "ObjectType": "network traffic",
                    "tactic": "Exfiltration",
                    "technique": [
                        "T1041-Exfiltration Over C2 Channel"
                    ],
                    "id": 26
                }
            ]
        }
    }
}