# cti_preprocessing/original_reports/ANALYSIS-ON-APT-TO-BE-ATTACK-THAT-FOCUSING-ON-CHINAS-GOVERNMENT-AGENCY-.json
{
    "triplets": [
        {
            "Subject": "APT-TOCS",
            "SubjectType": "threat-actor",
            "Relation": "use",
            "Object": "powershell.exe",
            "ObjectType": "tool",
            "tactic": "Initial Access",
            "technique": [
                "T1190-Exploit Public-Facing Application"
            ]
        },
        {
            "Subject": "Attackers",
            "SubjectType": "threat-actor",
            "Relation": "exploit",
            "Object": "security vulnerabilities",
            "ObjectType": "vulnerability",
            "tactic": "Initial Access",
            "technique": [
                "T1190-Exploit Public-Facing Application"
            ]
        },
        {
            "Subject": "Sample A",
            "SubjectType": "file",
            "Relation": "identified-as",
            "Object": "Trojan/Win32.MSShell",
            "ObjectType": "Malware",
            "tactic": "Initial Access",
            "technique": [
                "Others"
            ]
        },
        {
            "Subject": "Sample A",
            "SubjectType": "file",
            "Relation": "use",
            "Object": "WinExec",
            "ObjectType": "tool",
            "tactic": "Initial Access",
            "technique": [
                "T1190-Exploit Public-Facing Application"
            ]
        },
        {
            "Subject": "APT-TOCS attack",
            "SubjectType": "campaign",
            "Relation": "involve",
            "Object": "Shellcode script functionalities",
            "ObjectType": "artifacts",
            "tactic": "Execution",
            "technique": [
                "T1059-Command and Scripting Interpreter",
                "T1559-Inter-Process Communication"
            ]
        },
        {
            "Subject": "Shellcode script functionalities",
            "SubjectType": "artifacts",
            "Relation": "execute",
            "Object": "in memory",
            "ObjectType": "Location",
            "tactic": "Execution",
            "technique": [
                "T1059-Command and Scripting Interpreter",
                "T1106-Native API"
            ]
        },
        {
            "Subject": "attack pattern",
            "SubjectType": "attack-pattern",
            "Relation": "communicate-with",
            "Object": "Beacon mode",
            "ObjectType": "network traffic",
            "tactic": "Execution",
            "technique": [
                "T1053-Scheduled Task/Job"
            ]
        },
        {
            "Subject": "APT-TOCS",
            "SubjectType": "Malware",
            "Relation": "exhibit",
            "Object": "features",
            "ObjectType": "attack-pattern",
            "tactic": "Defense Evasion",
            "technique": [
                "T1564-Hide Artifacts"
            ]
        },
        {
            "Subject": "APT-TOCS",
            "SubjectType": "threat-actor",
            "Relation": "depend on",
            "Object": "network projection",
            "ObjectType": "technique",
            "tactic": "Lateral Movement",
            "technique": [
                "T1021-Remote Services"
            ]
        },
        {
            "Subject": "APT-TOCS",
            "SubjectType": "threat-actor",
            "Relation": "move laterally",
            "Object": "internal network",
            "ObjectType": "infrastructure",
            "tactic": "Lateral Movement",
            "technique": [
                "T1021-Remote Services"
            ]
        },
        {
            "Subject": "APT-TOCS",
            "SubjectType": "threat-actor",
            "Relation": "use",
            "Object": "encrypted data scripts",
            "ObjectType": "artifacts",
            "tactic": "Command and Control",
            "technique": [
                "T1071-Application Layer Protocol",
                "T1573-Encrypted Channel",
                "T1132-Data Encoding"
            ]
        },
        {
            "Subject": "PowerShell",
            "SubjectType": "tool",
            "Relation": "download",
            "Object": "encrypted data scripts",
            "ObjectType": "artifacts",
            "tactic": "Command and Control",
            "technique": [
                "T1071-Application Layer Protocol",
                "T1573-Encrypted Channel",
                "T1132-Data Encoding"
            ]
        },
        {
            "Subject": "The attack",
            "SubjectType": "campaign",
            "Relation": "communicate-with",
            "Object": "compromised systems",
            "ObjectType": "infrastrcucture",
            "tactic": "Command and Control",
            "technique": [
                "T1071-Application Layer Protocol"
            ]
        },
        {
            "Subject": "APT-TOCS",
            "SubjectType": "threat-actor",
            "Relation": "related-to",
            "Object": "Cobalt Strike",
            "ObjectType": "tool",
            "tactic": "Others",
            "technique": [
                "Others"
            ]
        },
        {
            "Subject": "Cobalt Strike",
            "SubjectType": "tool",
            "Relation": "penetrate",
            "Object": "firewalls",
            "ObjectType": "infrastructure",
            "tactic": "Others",
            "technique": [
                "Others"
            ]
        },
        {
            "Subject": "Cobalt Strike",
            "SubjectType": "tool",
            "Relation": "evade",
            "Object": "detection",
            "ObjectType": "course-of-action",
            "tactic": "Others",
            "technique": [
                "Others"
            ]
        }
    ],
    "file_name": "ANALYSIS-ON-APT-TO-BE-ATTACK-THAT-FOCUSING-ON-CHINAS-GOVERNMENT-AGENCY-",
    "tactic_label": true,
    "rewrite": {
        "Reconnaissance": "None",
        "Resource Development": "None",
        "Initial Access": "The APT-TOCS attack utilized \"powershell.exe\" to execute Shellcode scripts for remote control of targeted systems. Attackers may have used security vulnerabilities to execute scripts directly on the host. A binary file named Sample A, identified as Trojan/Win32.MSShell with MD5 44BCF2DD262F12222ADEAB6F59B2975B, was used in similar attacks to lead the attack by executing embedded malware using WinExec to call \"powershell.exe\".",
        "Execution": "The APT-TOCS attack's core step involved downloading Shellcode script functionalities, which were then executed in memory without writing to the hard disk. The attack pattern included communication via Beacon mode, sending network heartbeat packages every 60 seconds, and transmitting data information via the Cookie field.",
        "Persistence": "None",
        "Privilege Escalation": "None",
        "Defense Evasion": "APT-TOCS exhibited features designed to evade detection by security software and firewall interception. The attack pattern included running in memory, no hard disk writing operations, and the use of multi-beacon communication, which could work simultaneously and evade sandbox detection.",
        "Credential Access": "None",
        "Discovery": "None",
        "Lateral Movement": "APT-TOCS depended on network projection and could laterally move within the internal network as required to conduct an attack.",
        "Collection": "None",
        "Command and Control": "APT-TOCS used encrypted data scripts downloaded by PowerShell for command and control. The attack communicated with compromised systems using HTTP GET requests to specific URLs and transmitted heartbeat packages using Cookie fields with 60 seconds intervals.",
        "Exfiltration": "None",
        "Impact": "None",
        "Others": "APT-TOCS is closely related to the Cobalt Strike automatic attack testing platform. The attack can penetrate firewalls and is undetectable, capable of attacking various platforms, and formidable against Trusted Computing, Cloud Detection, and Sandbox Detection. The threat has been active for 5 years without powerful detection methods to defeat it. The attack is considered an APT due to its targeted nature, anti-detection functions, and ability to conceal itself. The use of a commercial attack platform reduces the cost and complexity of launching such attacks, making it difficult to trace the attackers."
    },
    "technique_label": true
}