# cti_preprocessing/original_reports/CopyKittens.json
{
    "triplets": [
        {
            "Subject": "CopyKittens attack group",
            "SubjectType": "threat-actor",
            "Relation": "initiated",
            "Object": "their cyber-attacks",
            "ObjectType": "campaign",
            "tactic": "Initial Access",
            "technique": [
                "T1566-Phishing"
            ]
        },
        {
            "Subject": "Spear phishing emails",
            "SubjectType": "Email Message",
            "Relation": "contain",
            "Object": "infected document files",
            "ObjectType": "file",
            "tactic": "Initial Access",
            "technique": [
                "T1566-Phishing"
            ]
        },
        {
            "Subject": "The documents",
            "SubjectType": "file",
            "Relation": "contained",
            "Object": "OLE binary objects",
            "ObjectType": "artifacts",
            "tactic": "Initial Access",
            "technique": [
                "T1566-Phishing"
            ]
        },
        {
            "Subject": "CopyKittens attack group",
            "SubjectType": "threat-actor",
            "Relation": "develop",
            "Object": "'homemade' tools",
            "ObjectType": "tool",
            "tactic": "Others",
            "technique": [
                "Others"
            ]
        },
        {
            "Subject": "CopyKittens attack group",
            "SubjectType": "threat-actor",
            "Relation": "target",
            "Object": "high-ranking diplomats and researchers specializing in Middle East Studies",
            "ObjectType": "identity",
            "tactic": "Others",
            "technique": [
                "Others"
            ]
        },
        {
            "Subject": "Matryoshka",
            "SubjectType": "tool",
            "Relation": "use by",
            "Object": "CopyKittens attack group",
            "ObjectType": "threat-actor",
            "tactic": "Others",
            "technique": [
                "Others"
            ]
        },
        {
            "Subject": "CopyKittens group",
            "SubjectType": "threat-actor",
            "Relation": "employ",
            "Object": "Matryoshka",
            "ObjectType": "campaign",
            "tactic": "Execution",
            "technique": [
                "Others"
            ]
        },
        {
            "Subject": "SCR PE file",
            "SubjectType": "file",
            "Relation": "signal-to",
            "Object": "C2 server",
            "ObjectType": "infrastructure",
            "tactic": "Execution",
            "technique": [
                "Others"
            ]
        },
        {
            "Subject": "Reflective Loader",
            "SubjectType": "tool",
            "Relation": "use",
            "Object": "anti-debugging and anti-sandboxing techniques",
            "ObjectType": "technique",
            "tactic": "Execution",
            "technique": [
                "Others"
            ]
        },
        {
            "Subject": "RAT component",
            "SubjectType": "tool",
            "Relation": "configure",
            "Object": "Reflective Loader",
            "ObjectType": "tool",
            "tactic": "Execution",
            "technique": [
                "Others"
            ]
        },
        {
            "Subject": "The RAT component",
            "SubjectType": "Malware",
            "Relation": "ensure",
            "Object": "persistence",
            "ObjectType": "Tactic",
            "tactic": "Persistence",
            "technique": [
                "T1547-Boot or Logon Autostart Execution",
                "T1053-Scheduled Task/Job"
            ]
        },
        {
            "Subject": "The RAT component",
            "SubjectType": "Malware",
            "Relation": "copy",
            "Object": "reflective loader",
            "ObjectType": "file",
            "tactic": "Persistence",
            "technique": [
                "T1547-Boot or Logon Autostart Execution"
            ]
        },
        {
            "Subject": "The RAT component",
            "SubjectType": "Malware",
            "Relation": "create",
            "Object": "registry key",
            "ObjectType": "registry",
            "tactic": "Persistence",
            "technique": [
                "T1547-Boot or Logon Autostart Execution"
            ]
        },
        {
            "Subject": "The RAT component",
            "SubjectType": "Malware",
            "Relation": "schedule",
            "Object": "Microsoft Boost Kernel Optimization",
            "ObjectType": "task",
            "tactic": "Persistence",
            "technique": [
                "T1053-Scheduled Task/Job"
            ]
        },
        {
            "Subject": "The CopyKittens group",
            "SubjectType": "threat-actor",
            "Relation": "use",
            "Object": "anti-debugging and anti-sandboxing methods",
            "ObjectType": "technique",
            "tactic": "Defense Evasion",
            "technique": [
                "T1497-Virtualization/Sandbox Evasion",
                "T1622-Debugger Evasion"
            ]
        },
        {
            "Subject": "The CopyKittens group",
            "SubjectType": "threat-actor",
            "Relation": "use",
            "Object": "runtime API address resolution",
            "ObjectType": "technique",
            "tactic": "Defense Evasion",
            "technique": [
                "T1497-Virtualization/Sandbox Evasion"
            ]
        },
        {
            "Subject": "The CopyKittens group",
            "SubjectType": "threat-actor",
            "Relation": "use",
            "Object": "global mutex",
            "ObjectType": "tool",
            "tactic": "Defense Evasion",
            "technique": [
                "T1564-Hide Artifacts"
            ]
        },
        {
            "Subject": "The RAT",
            "SubjectType": "Malware",
            "Relation": "capable of",
            "Object": "stealing Outlook passwords",
            "ObjectType": "attack-pattern",
            "tactic": "Credential Access",
            "technique": [
                "T1555-Credentials from Password Stores"
            ]
        },
        {
            "Subject": "The RAT",
            "SubjectType": "Malware",
            "Relation": "have",
            "Object": "screen grabbing and keylogging capabilities",
            "ObjectType": "attack-pattern",
            "tactic": "Collection",
            "technique": [
                "T1113-Screen Capture",
                "T1056-Input Capture"
            ]
        },
        {
            "Subject": "portions of the source code",
            "SubjectType": "artifacts",
            "Relation": "attributed-to",
            "Object": "rohitab.com",
            "ObjectType": "domain-name",
            "tactic": "Collection",
            "technique": [
                "Others"
            ]
        },
        {
            "Subject": "The RAT",
            "SubjectType": "Malware",
            "Relation": "use",
            "Object": "DNS protocol",
            "ObjectType": "network traffic",
            "tactic": "Command and Control",
            "technique": [
                "T1001-Data Obfuscation"
            ]
        },
        {
            "Subject": "The RAT",
            "SubjectType": "Malware",
            "Relation": "communicate-with",
            "Object": "C2 server",
            "ObjectType": "infrastructure",
            "tactic": "Command and Control",
            "technique": [
                "T1001-Data Obfuscation"
            ]
        },
        {
            "Subject": "The group",
            "SubjectType": "threat-actor",
            "Relation": "use",
            "Object": "substitute ciphers",
            "ObjectType": "technique",
            "tactic": "Command and Control",
            "technique": [
                "T1001-Data Obfuscation"
            ]
        },
        {
            "Subject": "RAT",
            "SubjectType": "Malware",
            "Relation": "exfiltrate-to",
            "Object": "DNS protocol",
            "ObjectType": "network traffic",
            "tactic": "Exfiltration",
            "technique": [
                "T1048-Exfiltration Over Alternative Protocol"
            ]
        }
    ],
    "file_name": "CopyKittens",
    "tactic_label": true,
    "rewrite": {
        "Reconnaissance": "None",
        "Resource Development": "None",
        "Initial Access": "The CopyKittens attack group initiated their cyber-attacks with spear phishing emails containing infected document files as attachments. These emails were carefully crafted with subjects that matched the interests of the targets, such as high-ranking diplomats and academic researchers. The documents contained OLE binary objects with deceptive names using the Right-To-Left Override Unicode character to appear as harmless file types.",
        "Execution": "The CopyKittens group employed a multi-stage infection framework named 'Matryoshka'. The first stage involved a dropper, an SCR PE file, which obfuscated code and signaled to the C2 server upon execution. It launched a loader to execute functions and compared anti-analysis logic, reporting back to C2. The second stage, the Reflective Loader, used anti-debugging and anti-sandboxing techniques, runtime API address resolution, covert DLL injection of the RAT library, and ensured persistence on disk. The third stage was the RAT component, which configured the Reflective Loader to survive reboots and process exits, established DNS Command and Control communication, and provided common RAT functionalities.",
        "Persistence": "The RAT component of the Matryoshka framework ensured its persistence by copying the reflective loader to common Windows folders and creating a registry key for automatic execution upon system boot. Additionally, it scheduled a task named 'Microsoft Boost Kernel Optimization' to rerun the injection routine every 20 minutes, using a global mutex to prevent multiple instances from causing instability.",
        "Privilege Escalation": "None",
        "Defense Evasion": "The CopyKittens group utilized techniques to avoid detection, such as employing anti-debugging and anti-sandboxing methods, and runtime API address resolution with obfuscated strings. They also used a global mutex to reduce the risk of multiple RAT instances running simultaneously.",
        "Credential Access": "The RAT was capable of stealing Outlook passwords, a method that may have been copied from SecurityExploded's description for 'Recovering Passwords from Outlook 2002-2013'.",
        "Discovery": "None",
        "Lateral Movement": "None",
        "Collection": "The RAT had capabilities for screen grabbing and keylogging, with portions of the source code traced back to the rohitab.com online forum.",
        "Command and Control": "The RAT used the DNS protocol to communicate with the C2 server, constructing DNS queries with a unique ID of the infected machine, a random string, and obfuscated data to be transmitted. The group used substitute ciphers to disguise the DNS traffic and made it difficult for traffic analysis and detection by using IP addresses from Microsoft and McAfee in the C2 responses.",
        "Exfiltration": "Data exfiltration was performed over the DNS protocol, with the RAT translating received commands from the C2 server into corresponding actions, such as stealing Outlook passwords when a specific IP address was received in the DNS response.",
        "Impact": "None",
        "Others": "The CopyKittens attack group is a mid-level threat actor known for their effective and stealthy multi-stage attacks, avoiding the use of known RATs and packers by developing 'homemade' tools. They have shown constant development to overcome security product improvements. The group's targets have included high-ranking diplomats and researchers specializing in Middle East Studies. The group's attack tool, 'Matryoshka', is a multi-stage framework with each part built to integrate with the subsequent step. The group's activities suggest that they will continue to target similar profiles in the future."
    },
    "technique_label": true
}