# cti_preprocessing/original_reports/FTA 1007 - Shamoon.json
{
    "triplets": [
        {
            "Subject": "Shamoon",
            "SubjectType": "Malware",
            "Relation": "deployed against",
            "Object": "entities",
            "ObjectType": "identity",
            "tactic": "Initial Access",
            "technique": [
                "T1078-Valid Accounts"
            ]
        },
        {
            "Subject": "Shamoon",
            "SubjectType": "Malware",
            "Relation": "designed to",
            "Object": "overwrite data",
            "ObjectType": "artifacts",
            "tactic": "Initial Access",
            "technique": [
                "T1078-Valid Accounts"
            ]
        },
        {
            "Subject": "Shamoon",
            "SubjectType": "Malware",
            "Relation": "executed as",
            "Object": "scheduled job",
            "ObjectType": "process",
            "tactic": "Initial Access",
            "technique": [
                "T1078-Valid Accounts"
            ]
        },
        {
            "Subject": "Shamoon",
            "SubjectType": "Malware",
            "Relation": "entrenched itself as",
            "Object": "service",
            "ObjectType": "software",
            "tactic": "Initial Access",
            "technique": [
                "T1078-Valid Accounts"
            ]
        },
        {
            "Subject": "Shamoon",
            "SubjectType": "Malware",
            "Relation": "dropped",
            "Object": "driver",
            "ObjectType": "file",
            "tactic": "Initial Access",
            "technique": [
                "Others"
            ]
        },
        {
            "Subject": "Shamoon malware",
            "SubjectType": "Malware",
            "Relation": "execute",
            "Object": "itself",
            "ObjectType": "Malware",
            "tactic": "Execution",
            "technique": [
                "T1053-Scheduled Task/Job"
            ]
        },
        {
            "Subject": "Shamoon malware",
            "SubjectType": "Malware",
            "Relation": "entrench",
            "Object": "itself",
            "ObjectType": "Malware",
            "tactic": "Execution",
            "technique": [
                "T1053-Scheduled Task/Job"
            ]
        },
        {
            "Subject": "Shamoon malware",
            "SubjectType": "Malware",
            "Relation": "drop",
            "Object": "driver",
            "ObjectType": "file",
            "tactic": "Execution",
            "technique": [
                "T1053-Scheduled Task/Job"
            ]
        },
        {
            "Subject": "Shamoon",
            "SubjectType": "Malware",
            "Relation": "entrench",
            "Object": "service",
            "ObjectType": "artifact",
            "tactic": "Persistence",
            "technique": [
                "Others"
            ]
        },
        {
            "Subject": "Shamoon malware",
            "SubjectType": "Malware",
            "Relation": "overwrite",
            "Object": "Documents and Settings folder",
            "ObjectType": "directory",
            "tactic": "Defense Evasion",
            "technique": [
                "T1484-Domain Policy Modification"
            ]
        },
        {
            "Subject": "Shamoon malware",
            "SubjectType": "Malware",
            "Relation": "overwrite",
            "Object": "System32 folder",
            "ObjectType": "directory",
            "tactic": "Defense Evasion",
            "technique": [
                "T1484-Domain Policy Modification"
            ]
        },
        {
            "Subject": "Shamoon malware",
            "SubjectType": "Malware",
            "Relation": "overwrite",
            "Object": "Master Boot Record",
            "ObjectType": "artifacts",
            "tactic": "Defense Evasion",
            "technique": [
                "T1484-Domain Policy Modification"
            ]
        },
        {
            "Subject": "Shamoon malware",
            "SubjectType": "Malware",
            "Relation": "overwrite",
            "Object": "Volume Boot Record",
            "ObjectType": "artifacts",
            "tactic": "Defense Evasion",
            "technique": [
                "T1006-Direct Volume Access"
            ]
        },
        {
            "Subject": "Fidelis researchers",
            "SubjectType": "author",
            "Relation": "develop",
            "Object": "rules",
            "ObjectType": "course-of-action",
            "tactic": "Command and Control",
            "technique": [
                "Others"
            ]
        },
        {
            "Subject": "rules",
            "SubjectType": "course-of-action",
            "Relation": "detect",
            "Object": "Shamoon malware",
            "ObjectType": "Malware",
            "tactic": "Command and Control",
            "technique": [
                "Others"
            ]
        },
        {
            "Subject": "malware detection engine",
            "SubjectType": "tool",
            "Relation": "recognize",
            "Object": "variant of Shamoon malware",
            "ObjectType": "Malware",
            "tactic": "Command and Control",
            "technique": [
                "Others"
            ]
        },
        {
            "Subject": "Shamoon",
            "SubjectType": "Malware",
            "Relation": "compromise",
            "Object": "file systems",
            "ObjectType": "artifact",
            "tactic": "Others",
            "technique": [
                "Others"
            ]
        },
        {
            "Subject": "Fidelis researchers",
            "SubjectType": "author",
            "Relation": "test",
            "Object": "strategies for VBR and file system recovery",
            "ObjectType": "course-of-action",
            "tactic": "Others",
            "technique": [
                "Others"
            ]
        },
        {
            "Subject": "Shamoon malware",
            "SubjectType": "Malware",
            "Relation": "overwrite",
            "Object": "data on disks",
            "ObjectType": "artifacts",
            "tactic": "Impact",
            "technique": [
                "T1485-Data Destruction"
            ]
        },
        {
            "Subject": "Shamoon malware",
            "SubjectType": "Malware",
            "Relation": "target",
            "Object": "user data",
            "ObjectType": "artifacts",
            "tactic": "Impact",
            "technique": [
                "T1485-Data Destruction"
            ]
        },
        {
            "Subject": "Shamoon malware",
            "SubjectType": "Malware",
            "Relation": "target",
            "Object": "system data",
            "ObjectType": "artifacts",
            "tactic": "Impact",
            "technique": [
                "T1485-Data Destruction"
            ]
        },
        {
            "Subject": "Shamoon malware",
            "SubjectType": "Malware",
            "Relation": "overwrite",
            "Object": "disk's boot records",
            "ObjectType": "artifacts",
            "tactic": "Impact",
            "technique": [
                "T1485-Data Destruction"
            ]
        },
        {
            "Subject": "Shamoon malware",
            "SubjectType": "Malware",
            "Relation": "prevent",
            "Object": "accessing the targeted file system",
            "ObjectType": "artifacts",
            "tactic": "Impact",
            "technique": [
                "T1485-Data Destruction"
            ]
        },
        {
            "Subject": "Shamoon malware",
            "SubjectType": "Malware",
            "Relation": "overwrite",
            "Object": "file data",
            "ObjectType": "artifacts",
            "tactic": "Impact",
            "technique": [
                "T1485-Data Destruction"
            ]
        }
    ],
    "file_name": "FTA 1007 - Shamoon",
    "tactic_label": true,
    "rewrite": {
        "Reconnaissance": "None",
        "Resource Development": "None",
        "Initial Access": "The Shamoon malware was deployed against certain entities around August 15, 2012. It was designed to overwrite data on disks attached to or accessible from targeted systems. The malware executed as a scheduled job, entrenched itself as a service, and dropped a driver that facilitated disk access to overwrite disk data, including the contents of the Documents and Settings and Windows\\system32 directories.",
        "Execution": "Shamoon malware executes a copy of itself as a scheduled job, deletes the file created for the job, entrenches itself as a service, and executes the entrenched file resulting in a dropped driver. The dropped driver is then loaded and executed, facilitating disk access for the malware's destructive operations.",
        "Persistence": "Shamoon entrenches itself as a service on the infected system, which is a method used to maintain persistence on the host.",
        "Privilege Escalation": "None",
        "Defense Evasion": "Shamoon malware's wiping function overwrites data within the Documents and Settings folder followed by the System32 folder, and then starts physical disk access to begin wiping at the disk level. The malware also overwrites the disk's boot records, including the Master Boot Record (MBR) and Volume Boot Record (VBR), to prevent access to the affected file systems.",
        "Credential Access": "None",
        "Discovery": "None",
        "Lateral Movement": "None",
        "Collection": "None",
        "Command and Control": "Fidelis researchers have developed rules for detecting Shamoon malware along the entire threat life cycle, including initial infection, lateral propagation, and command and control communication. The embedded malware detection engine recognizes the variant of Shamoon malware analyzed.",
        "Exfiltration": "None",
        "Impact": "Shamoon malware's impact includes overwriting data on disks, targeting user data first, then system data, and eventually overwriting the disk's boot records. This results in the prevention of accessing the targeted file system and substantial amounts of file data being overwritten.",
        "Others": "Recovery efforts for Shamoon-infected file systems were successful due to the malware's inability to completely wipe the disk before the infected host restarted. This allowed for the recovery of the majority of files outside of the System32 and Documents and Settings folders. Fidelis researchers tested several strategies for VBR and file system recovery, finding that the Master File Table (MFT) was mostly intact and that a backup copy of the VBR could be used to rebuild the file system. The recovery process is detailed, including the use of forensic software like EnCase and hex editors to restore the file system and recover relevant artifacts for forensic examination."
    },
    "technique_label": true
}