# cti_preprocessing/original_reports/greatidea_A41_v1.0.json
{
    "triplets": [
        {
            "Subject": "A41APT campaign",
            "SubjectType": "Campaign",
            "Relation": "utilize",
            "Object": "SSL-VPN abuse",
            "ObjectType": "Technique",
            "tactic": "Initial Access",
            "technique": [
                "T1133-External Remote Services",
                "T1078-Valid Accounts"
            ]
        },
        {
            "Subject": "A41APT campaign",
            "SubjectType": "Campaign",
            "Relation": "target",
            "Object": "Japanese companies",
            "ObjectType": "Identity",
            "tactic": "Initial Access",
            "technique": [
                "T1133-External Remote Services",
                "T1078-Valid Accounts"
            ]
        },
        {
            "Subject": "attackers",
            "SubjectType": "threat-actor",
            "Relation": "use",
            "Object": "DESLoader",
            "ObjectType": "Malware",
            "tactic": "Execution",
            "technique": [
                "T1129-Shared Modules"
            ]
        },
        {
            "Subject": "DESLoader",
            "SubjectType": "Malware",
            "Relation": "employ",
            "Object": "DLL side-loading techniques",
            "ObjectType": "Technique",
            "tactic": "Execution",
            "technique": [
                "T1129-Shared Modules"
            ]
        },
        {
            "Subject": "DESLoader",
            "SubjectType": "Malware",
            "Relation": "use",
            "Object": "fileless implants",
            "ObjectType": "Technique",
            "tactic": "Execution",
            "technique": [
                "T1129-Shared Modules"
            ]
        },
        {
            "Subject": "SodaMaster",
            "SubjectType": "Malware",
            "Relation": "execute",
            "Object": "commands",
            "ObjectType": "attack-pattern",
            "tactic": "Execution",
            "technique": [
                "T1059-Command and Scripting Interpreter"
            ]
        },
        {
            "Subject": "SodaMaster",
            "SubjectType": "Malware",
            "Relation": "communicate-with",
            "Object": "C2 servers",
            "ObjectType": "infrastructure",
            "tactic": "Execution",
            "technique": [
                "Others"
            ]
        },
        {
            "Subject": "attackers",
            "SubjectType": "threat-actor",
            "Relation": "register",
            "Object": "scheduled tasks",
            "ObjectType": "artifact",
            "tactic": "Persistence",
            "technique": [
                "T1053-Scheduled Task/Job"
            ]
        },
        {
            "Subject": "scheduled tasks",
            "SubjectType": "artifact",
            "Relation": "execute",
            "Object": "legitimate executable files",
            "ObjectType": "file",
            "tactic": "Persistence",
            "technique": [
                "T1053-Scheduled Task/Job"
            ]
        },
        {
            "Subject": "legitimate executable files",
            "SubjectType": "file",
            "Relation": "load",
            "Object": "DESLoader",
            "ObjectType": "Malware",
            "tactic": "Persistence",
            "technique": [
                "Others"
            ]
        },
        {
            "Subject": "scheduled tasks",
            "SubjectType": "artifact",
            "Relation": "mimic",
            "Object": "legitimate Windows tasks",
            "ObjectType": "artifact",
            "tactic": "Persistence",
            "technique": [
                "T1053-Scheduled Task/Job"
            ]
        },
        {
            "Subject": "A41APT campaign",
            "SubjectType": "campaign",
            "Relation": "exhibit",
            "Object": "defense evasion techniques",
            "ObjectType": "Tactic",
            "tactic": "Defense Evasion",
            "technique": [
                "T1027-Obfuscated Files or Information"
            ]
        },
        {
            "Subject": "attackers",
            "SubjectType": "threat-actor",
            "Relation": "use",
            "Object": "obfuscated PowerShell commands",
            "ObjectType": "tool",
            "tactic": "Defense Evasion",
            "technique": [
                "T1027-Obfuscated Files or Information"
            ]
        },
        {
            "Subject": "attackers",
            "SubjectType": "threat-actor",
            "Relation": "perform",
            "Object": "internal network scanning",
            "ObjectType": "attack-pattern",
            "tactic": "Credential Access",
            "technique": [
                "T1003-OS Credential Dumping"
            ]
        },
        {
            "Subject": "attackers",
            "SubjectType": "threat-actor",
            "Relation": "use",
            "Object": "administrator accounts",
            "ObjectType": "identity",
            "tactic": "Credential Access",
            "technique": [
                "T1003-OS Credential Dumping"
            ]
        },
        {
            "Subject": "attackers",
            "SubjectType": "threat-actor",
            "Relation": "deploy",
            "Object": "RDP",
            "ObjectType": "tool",
            "tactic": "Credential Access",
            "technique": [
                "T1003-OS Credential Dumping"
            ]
        },
        {
            "Subject": "attackers",
            "SubjectType": "threat-actor",
            "Relation": "execute",
            "Object": "csvde.exe",
            "ObjectType": "tool",
            "tactic": "Credential Access",
            "technique": [
                "T1003-OS Credential Dumping"
            ]
        },
        {
            "Subject": "attackers",
            "SubjectType": "threat-actor",
            "Relation": "execute",
            "Object": "AdFind",
            "ObjectType": "tool",
            "tactic": "Credential Access",
            "technique": [
                "T1003-OS Credential Dumping"
            ]
        },
        {
            "Subject": "attackers",
            "SubjectType": "threat-actor",
            "Relation": "dump",
            "Object": "credentials",
            "ObjectType": "artifacts",
            "tactic": "Credential Access",
            "technique": [
                "T1003-OS Credential Dumping"
            ]
        },
        {
            "Subject": "attackers",
            "SubjectType": "threat-actor",
            "Relation": "conduct",
            "Object": "discovery activities",
            "ObjectType": "attack-pattern",
            "tactic": "Discovery",
            "technique": [
                "T1046-Network Service Discovery"
            ]
        },
        {
            "Subject": "attackers",
            "SubjectType": "threat-actor",
            "Relation": "use",
            "Object": "WMIC",
            "ObjectType": "tool",
            "tactic": "Discovery",
            "technique": [
                "T1046-Network Service Discovery",
                "T1518-Software Discovery"
            ]
        },
        {
            "Subject": "attackers",
            "SubjectType": "threat-actor",
            "Relation": "deploy",
            "Object": "tools",
            "ObjectType": "tool",
            "tactic": "Discovery",
            "technique": [
                "T1046-Network Service Discovery",
                "T1482-Domain Trust Discovery"
            ]
        },
        {
            "Subject": "The A41APT campaign",
            "SubjectType": "campaign",
            "Relation": "utilized",
            "Object": "RDP",
            "ObjectType": "tool",
            "tactic": "Lateral Movement",
            "technique": [
                "T1021-Remote Services"
            ]
        },
        {
            "Subject": "The A41APT campaign",
            "SubjectType": "campaign",
            "Relation": "target",
            "Object": "servers",
            "ObjectType": "infrastructure",
            "tactic": "Lateral Movement",
            "technique": [
                "T1021-Remote Services"
            ]
        },
        {
            "Subject": "They",
            "SubjectType": "campaign",
            "Relation": "used",
            "Object": "network scanning",
            "ObjectType": "attack-pattern",
            "tactic": "Lateral Movement",
            "technique": [
                "Others"
            ]
        },
        {
            "Subject": "They",
            "SubjectType": "campaign",
            "Relation": "used",
            "Object": "credential theft",
            "ObjectType": "attack-pattern",
            "tactic": "Lateral Movement",
            "technique": [
                "T1550-Use Alternate Authentication Material"
            ]
        },
        {
            "Subject": "attackers",
            "SubjectType": "threat-actor",
            "Relation": "communicate-with",
            "Object": "C2 servers",
            "ObjectType": "infrastructure",
            "tactic": "Command and Control",
            "technique": [
                "Others"
            ]
        },
        {
            "Subject": "attackers",
            "SubjectType": "threat-actor",
            "Relation": "use",
            "Object": "IP addresses",
            "ObjectType": "ipv4-addr",
            "tactic": "Command and Control",
            "technique": [
                "Others"
            ]
        },
        {
            "Subject": "APT10",
            "SubjectType": "threat-actor",
            "Relation": "attributed-to",
            "Object": "A41APT campaign",
            "ObjectType": "campaign",
            "tactic": "Others",
            "technique": [
                "Others"
            ]
        },
        {
            "Subject": "BlackTech",
            "SubjectType": "threat-actor",
            "Relation": "relevance",
            "Object": "A41APT campaign",
            "ObjectType": "campaign",
            "tactic": "Others",
            "technique": [
                "Others"
            ]
        },
        {
            "Subject": "A41APT campaign",
            "SubjectType": "campaign",
            "Relation": "target",
            "Object": "industries",
            "ObjectType": "infrastrcucture",
            "tactic": "Others",
            "technique": [
                "Others"
            ]
        },
        {
            "Subject": "A41APT campaign",
            "SubjectType": "campaign",
            "Relation": "use",
            "Object": "IP addresses",
            "ObjectType": "ipv4-addr",
            "tactic": "Others",
            "technique": [
                "Others"
            ]
        }
    ],
    "file_name": "greatidea_A41_v1.0",
    "tactic_label": true,
    "rewrite": {
        "Reconnaissance": "None",
        "Resource Development": "None",
        "Initial Access": "The A41APT campaign utilized SSL-VPN abuse to gain initial access to the target network, specifically targeting Japanese companies including their overseas branches. The attackers did not rely on spear-phishing but instead exploited SSL-VPN vulnerabilities or used stolen credentials to penetrate the network. The hostname 'DESKTOP-A41UVJV' was consistently used during these initial intrusions.",
        "Execution": "The attackers executed malicious code through the use of DESLoader, which is a multi-layer loader for various payloads including SodaMaster, P8RAT, and FYAntiLoader. DESLoader employed DLL side-loading techniques and fileless implants to execute its payloads. The payloads themselves, such as SodaMaster, were fileless RATs that executed commands and communicated with C2 servers.",
        "Persistence": "To maintain their foothold, the attackers registered scheduled tasks that executed legitimate executable files to load DESLoader every 15 minutes. The scheduled tasks were created with names that mimicked legitimate Windows tasks, making them harder to detect.",
        "Privilege Escalation": "None",
        "Defense Evasion": "The A41APT campaign exhibited sophisticated defense evasion techniques, including the use of fileless implants, DLL side-loading, and the removal of event logs after C2 communication. The attackers also used obfuscated PowerShell commands to delete event logs, further covering their tracks.",
        "Credential Access": "The attackers performed internal network scanning to find open RDP and SMB ports after initial access via SSL-VPN. They used administrator accounts to deploy RDP to servers and executed tools like csvde.exe and AdFind to dump credentials from systems such as AD servers, file servers, and backup servers.",
        "Discovery": "The attackers conducted discovery activities post-compromise, including scanning for open RDP and SMB ports, using WMIC to collect information on security products, and deploying tools to discover domain accounts and trust relationships.",
        "Lateral Movement": "The A41APT campaign heavily utilized RDP for lateral movement within the network, mainly targeting servers. They used network scanning and credential theft to move laterally and gain access to additional systems.",
        "Collection": "None",
        "Command and Control": "The attackers communicated with C2 servers via DESLoader's payloads or PowerShell remoting. They favored the use of IP addresses over domains for C2 infrastructure and showed a tendency not to reuse IP addresses repeatedly.",
        "Exfiltration": "None",
        "Impact": "None",
        "Others": "The A41APT campaign is attributed with strong associations to APT10 and potential relevance to BlackTech. The campaign targeted a wide range of industries and exhibited capabilities such as abusing DLL side-loading and removing traces post-compromise. The attackers used heavy IP address usage for C2 and avoided domain usage, with little reuse of IP addresses for C2. The IP used for initial intrusion and C2 were different. Countermeasures against this campaign include the implementation of MFA, patch adaptation operations, network monitoring, and strengthening security measures for servers."
    },
    "technique_label": true
}