# cti_preprocessing/original_reports/Iran’s APT34 Returns with an Updated Arsenal - Check Point Research.json
{
    "triplets": [
        {
            "Subject": "APT34",
            "SubjectType": "threat-actor",
            "Relation": "update",
            "Object": "payload arsenal",
            "ObjectType": "tool",
            "tactic": "Others",
            "technique": [
                "Others"
            ]
        },
        {
            "Subject": "APT34",
            "SubjectType": "threat-actor",
            "Relation": "target",
            "Object": "individuals",
            "ObjectType": "identity",
            "tactic": "Others",
            "technique": [
                "Others"
            ]
        },
        {
            "Subject": "SideTwist",
            "SubjectType": "Malware",
            "Relation": "provide",
            "Object": "functionality",
            "ObjectType": "tool",
            "tactic": "Others",
            "technique": [
                "Others"
            ]
        },
        {
            "Subject": "Karkoff",
            "SubjectType": "Malware",
            "Relation": "use",
            "Object": "exchange servers",
            "ObjectType": "infrastructure",
            "tactic": "Others",
            "technique": [
                "Others"
            ]
        },
        {
            "Subject": "Lebanese Government",
            "SubjectType": "identity",
            "Relation": "own",
            "Object": "mail.army.gov[.]lb",
            "ObjectType": "domain-name",
            "tactic": "Others",
            "technique": [
                "Others"
            ]
        },
        {
            "Subject": "Check Point Sandblast",
            "SubjectType": "CTI vender",
            "Relation": "provide",
            "Object": "protection",
            "ObjectType": "course-of-action",
            "tactic": "Others",
            "technique": [
                "Others"
            ]
        },
        {
            "Subject": "Check Point Research",
            "SubjectType": "CTI vender",
            "Relation": "identified",
            "Object": "SideTwist",
            "ObjectType": "Malware",
            "tactic": "Initial Access",
            "technique": [
                "T1566-Phishing"
            ]
        },
        {
            "Subject": "APT34",
            "SubjectType": "threat-actor",
            "Relation": "target",
            "Object": "a Lebanese entity",
            "ObjectType": "identity",
            "tactic": "Initial Access",
            "technique": [
                "T1566-Phishing"
            ]
        },
        {
            "Subject": "APT34",
            "SubjectType": "threat-actor",
            "Relation": "use",
            "Object": "Job-Details.doc",
            "ObjectType": "file",
            "tactic": "Initial Access",
            "technique": [
                "T1566-Phishing"
            ]
        },
        {
            "Subject": "Job-Details.doc",
            "SubjectType": "file",
            "Relation": "deliver",
            "Object": "SideTwist",
            "ObjectType": "Malware",
            "tactic": "Initial Access",
            "technique": [
                "T1566-Phishing"
            ]
        },
        {
            "Subject": "initial stage macros",
            "SubjectType": "artifact",
            "Relation": "establish",
            "Object": "SystemFailureReporter",
            "ObjectType": "artifact",
            "tactic": "Persistence",
            "technique": [
                "T1053-Scheduled Task/Job"
            ]
        },
        {
            "Subject": "SystemFailureReporter",
            "SubjectType": "artifact",
            "Relation": "execute",
            "Object": "SideTwist",
            "ObjectType": "Malware",
            "tactic": "Persistence",
            "technique": [
                "T1053-Scheduled Task/Job"
            ]
        },
        {
            "Subject": "APT34's campaign",
            "SubjectType": "campaign",
            "Relation": "utilize",
            "Object": "DNS tunneling",
            "ObjectType": "technique",
            "tactic": "Defense Evasion",
            "technique": [
                "T1205-Traffic Signaling"
            ]
        },
        {
            "Subject": "initial macros stage",
            "SubjectType": "attack-pattern",
            "Relation": "use",
            "Object": "DNS requests",
            "ObjectType": "network traffic",
            "tactic": "Defense Evasion",
            "technique": [
                "T1205-Traffic Signaling"
            ]
        },
        {
            "Subject": "DNS requests",
            "SubjectType": "network trafic",
            "Relation": "beacon-to",
            "Object": "attacker",
            "ObjectType": "threat-actor",
            "tactic": "Defense Evasion",
            "technique": [
                "T1205-Traffic Signaling"
            ]
        },
        {
            "Subject": "attacker",
            "SubjectType": "threat-actor",
            "Relation": "employ",
            "Object": "requestbin.net DNS tunneling service",
            "ObjectType": "tool",
            "tactic": "Defense Evasion",
            "technique": [
                "T1205-Traffic Signaling"
            ]
        },
        {
            "Subject": "The SideTwist backdoor",
            "SubjectType": "Malware",
            "Relation": "communicate-with",
            "Object": "the C&C server (sarmsoftware[.]com)",
            "ObjectType": "domain-name",
            "tactic": "Command and Control",
            "technique": [
                "T1071-Application Layer Protocol",
                "T1573-Encrypted Channel"
            ]
        },
        {
            "Subject": "malicious macros",
            "SubjectType": "artifact",
            "Relation": "execute",
            "Object": "series of actions",
            "ObjectType": "attack-pattern",
            "tactic": "Execution",
            "technique": [
                "T1059-Command and Scripting Interpreter"
            ]
        }
    ],
    "file_name": "Iran\u2019s APT34 Returns with an Updated Arsenal - Check Point Research",
    "tactic_label": true,
    "rewrite": {
        "Reconnaissance": "None",
        "Resource Development": "None",
        "Initial Access": "Check Point Research identified a new campaign by the Iranian threat group APT34 targeting a Lebanese entity using a novel backdoor variant named SideTwist. The initial access was achieved through a malicious Microsoft Word document titled 'Job-Details.doc' (md5: 6615c410b8d7411ed14946635947325e) delivered via LinkedIn messages. The document contained embedded malicious macros that, when activated, initiated the infection process.",
        "Execution": "The malicious macros within the Word document executed a series of actions upon activation: verifying a mouse connection (anti-sandboxing), fingerprinting the target device, communicating with the C2 server, dropping an embedded executable with a '.doc' extension (later renamed to '.exe'), and registering a Windows scheduled task to launch the executable periodically.",
        "Persistence": "Persistence was established by the initial stage macros through a scheduled task named 'SystemFailureReporter,' which was set to execute the second stage payload, SideTwist, every 5 minutes.",
        "Privilege Escalation": "None",
        "Defense Evasion": "APT34's campaign utilized DNS tunneling to avoid detection, with the initial macros stage using DNS requests to beacon back to the attacker and deliver victim identifiable information. The attacker employed the requestbin.net DNS tunneling service to receive updates about the infection progress without exposing their infrastructure.",
        "Credential Access": "None",
        "Discovery": "None",
        "Lateral Movement": "None",
        "Collection": "None",
        "Command and Control": "The SideTwist backdoor communicated with the C&C server (sarmsoftware[.]com) using HTTP on port 443, with port 80 as a fallback. The backdoor used different techniques for outgoing and incoming communications, both encrypted using the Mersenne Twister pseudorandom number generator. Commands were received from the server by searching for specific patterns hidden inside the HTML content of a legitimate-looking website, in this case, a fake Flickr page.",
        "Exfiltration": "None",
        "Impact": "None",
        "Others": "APT34, also known as OilRig, has been actively updating their payload arsenal post-2019 tool leak by Lab Dookhtegan. The group has a history of targeting individuals with job opportunity documents and has continued this method in their latest campaign. The SideTwist backdoor variant is new but provides functionality similar to other C-based backdoors used by APT34, such as DNSpionage and TONEDEAF. The campaign also involved another APT34 tool, a .NET-based backdoor named Karkoff, which used internet-facing exchange servers for communication. The Lebanese Government's exchange server (mail.army.gov[.]lb) credentials were found within the Karkoff implant, suggesting a potential long-term compromise. Check Point Sandblast provides protection against this APT attack."
    },
    "technique_label": true
}