# cti_preprocessing/original_reports/TA505 is Expanding its Operations.json
{
    "triplets": [
        {
            "Subject": "macro code",
            "SubjectType": "tool",
            "Relation": "download",
            "Object": "rtegre.exe and wprgxyeqd79.exe",
            "ObjectType": "file",
            "tactic": "Execution",
            "technique": [
                "T1204-User Execution"
            ]
        },
        {
            "Subject": "rtegre.exe and wprgxyeqd79.exe",
            "SubjectType": "file",
            "Relation": "store in",
            "Object": "C:\\Users\\Public directory",
            "ObjectType": "directory",
            "tactic": "Execution",
            "technique": [
                "T1204-User Execution"
            ]
        },
        {
            "Subject": "SFX",
            "SubjectType": "file",
            "Relation": "contain",
            "Object": "RMS client",
            "ObjectType": "tool",
            "tactic": "Execution",
            "technique": [
                "T1204-User Execution"
            ]
        },
        {
            "Subject": "i.cmd",
            "SubjectType": "file",
            "Relation": "execute",
            "Object": "uninstall.exe",
            "ObjectType": "file",
            "tactic": "Execution",
            "technique": [
                "T1204-User Execution"
            ]
        },
        {
            "Subject": "attack",
            "SubjectType": "attack-pattern",
            "Relation": "begin with",
            "Object": "spear phishing email",
            "ObjectType": "Email Message",
            "tactic": "Initial Access",
            "technique": [
                "T1566-Phishing"
            ]
        },
        {
            "Subject": "spear phishing email",
            "SubjectType": "Email Message",
            "Relation": "contain",
            "Object": "spreadsheet",
            "ObjectType": "file",
            "tactic": "Initial Access",
            "technique": [
                "T1566-Phishing"
            ]
        },
        {
            "Subject": "spreadsheet",
            "SubjectType": "file",
            "Relation": "contain",
            "Object": "malicious macro",
            "ObjectType": "software",
            "tactic": "Initial Access",
            "technique": [
                "T1566-Phishing"
            ]
        },
        {
            "Subject": "user",
            "SubjectType": "identity",
            "Relation": "attempt to view",
            "Object": "content",
            "ObjectType": "artifacts",
            "tactic": "Initial Access",
            "technique": [
                "T1566-Phishing"
            ]
        },
        {
            "Subject": "malicious macro",
            "SubjectType": "software",
            "Relation": "trigger",
            "Object": "infection process",
            "ObjectType": "attack-pattern",
            "tactic": "Initial Access",
            "technique": [
                "T1566-Phishing"
            ]
        },
        {
            "Subject": "attack",
            "SubjectType": "campaign",
            "Relation": "ensure",
            "Object": "persistence",
            "ObjectType": "Tactic",
            "tactic": "Persistence",
            "technique": [
                "T1547-Boot or Logon Autostart Execution"
            ]
        },
        {
            "Subject": "Windows Anytime Upgrade",
            "SubjectType": "registry",
            "Relation": "point to",
            "Object": "winserv.exe",
            "ObjectType": "file",
            "tactic": "Persistence",
            "technique": [
                "T1547-Boot or Logon Autostart Execution"
            ]
        },
        {
            "Subject": "winserv.exe",
            "SubjectType": "file",
            "Relation": "store in",
            "Object": "%ALLUSERSPROFILE%\\Windows Anytime Upgrade",
            "ObjectType": "directory",
            "tactic": "Persistence",
            "technique": [
                "Others"
            ]
        },
        {
            "Subject": "RMS client",
            "SubjectType": "tool",
            "Relation": "use",
            "Object": "MPress PE compressor utility",
            "ObjectType": "tool",
            "tactic": "Defense Evasion",
            "technique": [
                "T1027-Obfuscated Files or Information"
            ]
        },
        {
            "Subject": "RMS client",
            "SubjectType": "tool",
            "Relation": "kill",
            "Object": "rundll32.exe process",
            "ObjectType": "process",
            "tactic": "Defense Evasion",
            "technique": [
                "T1564-Hide Artifacts"
            ]
        },
        {
            "Subject": "RMS client",
            "SubjectType": "tool",
            "Relation": "establish",
            "Object": "command and control server",
            "ObjectType": "infrastructure",
            "tactic": "Command and Control",
            "technique": [
                "Others"
            ]
        },
        {
            "Subject": "command and control server",
            "SubjectType": "infrastructure",
            "Relation": "host",
            "Object": "217.12.201.159",
            "ObjectType": "ipv4-addr",
            "tactic": "Command and Control",
            "technique": [
                "Others"
            ]
        },
        {
            "Subject": "217.12.201.159",
            "SubjectType": "ipv4-addr",
            "Relation": "part-of",
            "Object": "AS-21100",
            "ObjectType": "autonomous-system",
            "tactic": "Command and Control",
            "technique": [
                "Others"
            ]
        },
        {
            "Subject": "AS-21100",
            "SubjectType": "autonomous-system",
            "Relation": "operate-by",
            "Object": "ITL LLC",
            "ObjectType": "identity",
            "tactic": "Command and Control",
            "technique": [
                "Others"
            ]
        },
        {
            "Subject": "TA505 group",
            "SubjectType": "threat-actor",
            "Relation": "target",
            "Object": "Banking and Retail sectors",
            "ObjectType": "identity",
            "tactic": "Others",
            "technique": [
                "Others"
            ]
        },
        {
            "Subject": "TA505 group",
            "SubjectType": "threat-actor",
            "Relation": "use",
            "Object": "RMS software",
            "ObjectType": "tool",
            "tactic": "Others",
            "technique": [
                "Others"
            ]
        }
    ],
    "file_name": "TA505 is Expanding its Operations",
    "tactic_label": true,
    "rewrite": {
        "Reconnaissance": "None",
        "Resource Development": "None",
        "Initial Access": "The intercepted attack began with a spear phishing email containing a spreadsheet with a malicious macro. The macro was triggered when the user attempted to view the content, initiating the infection process.",
        "Execution": "The macro code, upon execution, downloaded two files from the domain kentona[.su, specifically \"rtegre.exe\" and \"wprgxyeqd79.exe\", and stored them in the C:\\Users\\Public directory. These files were identified as a Trojan/Downloader and a Self-Extracting Archive (SFX), respectively. The SFX contained four files that were extracted to the %TEMP% folder and executed a batch script named \"i.cmd\" which in turn executed \"uninstall.exe\", another SFX archive, containing the RMS (Remote Manipulator System) client by TektonIT.",
        "Persistence": "The attack ensured persistence on the victim's machine by setting the \"HKCU\\Software\\Microsoft\\Windows\\CurrentVersion\\Run\" registry key with a new entry named \"Windows Anytime Upgrade\" pointing to \"winserv.exe\" stored in the \"%ALLUSERSPROFILE%\\Windows Anytime Upgrade\" directory.",
        "Privilege Escalation": "None",
        "Defense Evasion": "The RMS client, \"winserv.exe\", was encrypted using the MPress PE compressor utility to avoid antivirus detection. The script also continuously killed every \"rundll32.exe\" process on the victim machine to generate noise and avoid detection.",
        "Credential Access": "None",
        "Discovery": "None",
        "Lateral Movement": "None",
        "Collection": "None",
        "Command and Control": "The RMS client established a connection with a command and control server hosted on a Bulgarian remote host 217.12.201.159, part of a Virtual Dedicated Server subnet of the AS-21100, operated by ITL LLC.",
        "Exfiltration": "None",
        "Impact": "None",
        "Others": "The TA505 group, known for targeting the Banking and Retail sectors, may be expanding their operations to other industries. This conclusion is based on the similarities between the current attack and a previous campaign against a US retail company, which also used the RMS software. The campaign appears to be ongoing, with slight changes observed in the dropped files."
    },
    "technique_label": true
}