# cti_preprocessing/original_reports/the-konni-apt-chronicle-tracing-their-intelligence-driven-attack-chain.json
{
    "triplets": [
        {
            "Subject": "Konni APT Group",
            "SubjectType": "threat-actor",
            "Relation": "initiate",
            "Object": "attack",
            "ObjectType": "campaign",
            "tactic": "Initial Access",
            "technique": [
                "T1566-Phishing"
            ]
        },
        {
            "Subject": "Konni APT Group",
            "SubjectType": "threat-actor",
            "Relation": "deliver",
            "Object": "ISO file",
            "ObjectType": "file",
            "tactic": "Initial Access",
            "technique": [
                "T1566-Phishing"
            ]
        },
        {
            "Subject": "ISO file",
            "SubjectType": "file",
            "Relation": "contain",
            "Object": "zip file with malicious scripts",
            "ObjectType": "file",
            "tactic": "Initial Access",
            "technique": [
                "T1566-Phishing"
            ]
        },
        {
            "Subject": "install.vbs",
            "SubjectType": "file",
            "Relation": "execute",
            "Object": "resolve.dns",
            "ObjectType": "file",
            "tactic": "Execution",
            "technique": [
                "T1059-Command and Scripting Interpreter"
            ]
        },
        {
            "Subject": "update.vbs",
            "SubjectType": "file",
            "Relation": "execute",
            "Object": "stopedge.bat",
            "ObjectType": "file",
            "tactic": "Execution",
            "technique": [
                "T1059-Command and Scripting Interpreter"
            ]
        },
        {
            "Subject": "activate.vbs",
            "SubjectType": "file",
            "Relation": "send",
            "Object": "HTTP request",
            "ObjectType": "network traffic",
            "tactic": "Execution",
            "technique": [
                "T1059-Command and Scripting Interpreter"
            ]
        },
        {
            "Subject": "unactivate.vbs",
            "SubjectType": "file",
            "Relation": "capable of",
            "Object": "exfiltration",
            "ObjectType": "network traffic",
            "tactic": "Execution",
            "technique": [
                "T1059-Command and Scripting Interpreter"
            ]
        },
        {
            "Subject": "stopedge.bat",
            "SubjectType": "file",
            "Relation": "schedule",
            "Object": "install.vbs",
            "ObjectType": "file",
            "tactic": "Execution",
            "technique": [
                "T1059-Command and Scripting Interpreter",
                "T1053-Scheduled Task/Job"
            ]
        },
        {
            "Subject": "resolvedns.bat",
            "SubjectType": "file",
            "Relation": "initiate",
            "Object": "download",
            "ObjectType": "network traffic",
            "tactic": "Execution",
            "technique": [
                "T1059-Command and Scripting Interpreter"
            ]
        },
        {
            "Subject": "versioninfo.bat",
            "SubjectType": "file",
            "Relation": "collect",
            "Object": "information",
            "ObjectType": "artifact",
            "tactic": "Execution",
            "technique": [
                "T1059-Command and Scripting Interpreter"
            ]
        },
        {
            "Subject": "versioninfo.bat",
            "SubjectType": "file",
            "Relation": "exfiltrate",
            "Object": "data",
            "ObjectType": "artifact",
            "tactic": "Execution",
            "technique": [
                "T1059-Command and Scripting Interpreter"
            ]
        },
        {
            "Subject": "scripts",
            "SubjectType": "artifact",
            "Relation": "use",
            "Object": "XOR",
            "ObjectType": "technique",
            "tactic": "Defense Evasion",
            "technique": [
                "T1027-Obfuscated Files or Information",
                "T1070-Indicator Removal"
            ]
        },
        {
            "Subject": "scripts",
            "SubjectType": "artifact",
            "Relation": "use",
            "Object": "obfuscation",
            "ObjectType": "technique",
            "tactic": "Defense Evasion",
            "technique": [
                "T1027-Obfuscated Files or Information"
            ]
        },
        {
            "Subject": "stopedge.bat",
            "SubjectType": "file",
            "Relation": "manage",
            "Object": "deletion",
            "ObjectType": "technique",
            "tactic": "Defense Evasion",
            "technique": [
                "T1070-Indicator Removal"
            ]
        },
        {
            "Subject": "stopedge.bat",
            "SubjectType": "file",
            "Relation": "delete",
            "Object": "versioninfo.bat",
            "ObjectType": "file",
            "tactic": "Defense Evasion",
            "technique": [
                "T1070-Indicator Removal"
            ]
        },
        {
            "Subject": "stopedge.bat",
            "SubjectType": "file",
            "Relation": "delete",
            "Object": "update.vbs",
            "ObjectType": "file",
            "tactic": "Defense Evasion",
            "technique": [
                "T1070-Indicator Removal"
            ]
        },
        {
            "Subject": "'activate.vbs' script",
            "SubjectType": "file",
            "Relation": "communicate-with",
            "Object": "C2 server",
            "ObjectType": "infrastructure",
            "tactic": "Command and Control",
            "technique": [
                "T1071-Application Layer Protocol"
            ]
        },
        {
            "Subject": "'resolvedns.bat'",
            "SubjectType": "file",
            "Relation": "facilitate communication with",
            "Object": "C2 server",
            "ObjectType": "infrastructure",
            "tactic": "Command and Control",
            "technique": [
                "T1071-Application Layer Protocol"
            ]
        },
        {
            "Subject": "Konni APT Group",
            "SubjectType": "threat-actor",
            "Relation": "originate-from",
            "Object": "North Korea",
            "ObjectType": "Location",
            "tactic": "Others",
            "technique": [
                "Others"
            ]
        },
        {
            "Subject": "Konni APT Group",
            "SubjectType": "threat-actor",
            "Relation": "target",
            "Object": "South Korea",
            "ObjectType": "Location",
            "tactic": "Others",
            "technique": [
                "Others"
            ]
        },
        {
            "Subject": "stopedge.bat",
            "SubjectType": "file",
            "Relation": "contribute to",
            "Object": "persistence",
            "ObjectType": "Tactic",
            "tactic": "Persistence",
            "technique": [
                "T1053-Scheduled Task/Job"
            ]
        },
        {
            "Subject": "stopedge.bat",
            "SubjectType": "file",
            "Relation": "schedule",
            "Object": "install.vbs",
            "ObjectType": "file",
            "tactic": "Persistence",
            "technique": [
                "T1053-Scheduled Task/Job"
            ]
        },
        {
            "Subject": "install.vbs",
            "SubjectType": "file",
            "Relation": "run",
            "Object": "every 33 minutes",
            "ObjectType": "Technique",
            "tactic": "Persistence",
            "technique": [
                "T1053-Scheduled Task/Job"
            ]
        },
        {
            "Subject": "paycom.ini",
            "SubjectType": "file",
            "Relation": "present in",
            "Object": "directory",
            "ObjectType": "directory",
            "tactic": "Persistence",
            "technique": [
                "Others"
            ]
        },
        {
            "Subject": "versioninfo.bat",
            "SubjectType": "file",
            "Relation": "discover",
            "Object": "information about the user's computer",
            "ObjectType": "artifacts",
            "tactic": "Discovery",
            "technique": [
                "T1083-File and Directory Discovery"
            ]
        },
        {
            "Subject": "versioninfo.bat",
            "SubjectType": "file",
            "Relation": "collect",
            "Object": "information",
            "ObjectType": "artifacts",
            "tactic": "Collection",
            "technique": [
                "T1074-Data Staged"
            ]
        },
        {
            "Subject": "versioninfo.bat",
            "SubjectType": "file",
            "Relation": "prepare",
            "Object": "information",
            "ObjectType": "artifacts",
            "tactic": "Collection",
            "technique": [
                "T1074-Data Staged"
            ]
        },
        {
            "Subject": "information",
            "SubjectType": "artifacts",
            "Relation": "exfiltrate-to",
            "Object": "user's computer",
            "ObjectType": "infrastrcucture",
            "tactic": "Collection",
            "technique": [
                "T1074-Data Staged"
            ]
        },
        {
            "Subject": "unactivate.vbs",
            "SubjectType": "file",
            "Relation": "exfiltrate-to",
            "Object": "POST request",
            "ObjectType": "network traffic",
            "tactic": "Exfiltration",
            "technique": [
                "T1567-Exfiltration Over Web Service"
            ]
        },
        {
            "Subject": "versioninfo.bat",
            "SubjectType": "file",
            "Relation": "use",
            "Object": "unactivate.vbs",
            "ObjectType": "file",
            "tactic": "Exfiltration",
            "technique": [
                "T1567-Exfiltration Over Web Service"
            ]
        },
        {
            "Subject": "versioninfo.bat",
            "SubjectType": "file",
            "Relation": "exfiltrate-to",
            "Object": "collected information about the user's computer",
            "ObjectType": "artifacts",
            "tactic": "Exfiltration",
            "technique": [
                "T1567-Exfiltration Over Web Service"
            ]
        }
    ],
    "file_name": "the-konni-apt-chronicle-tracing-their-intelligence-driven-attack-chain",
    "tactic_label": true,
    "rewrite": {
        "Reconnaissance": "None",
        "Resource Development": "None",
        "Initial Access": "The Konni APT Group initiated their attack by delivering an ISO file downloaded from an in-the-wild (ITW) URL. This ISO file contained a zip file with malicious scripts, marking the start of the attack chain.",
        "Execution": "The attack involved the execution of various VBScripts and batch files. The 'install.vbs' script, once deobfuscated, executed a batch script named 'resolve.dns'. The 'update.vbs' script, similar to 'install.vbs', executed 'stopedge.bat'. The 'activate.vbs' script sent an HTTP request and saved the retrieved content. The 'unactivate.vbs' script was capable of exfiltration through a POST request. The 'stopedge.bat' file, if 'paycom.ini' existed, scheduled 'install.vbs' to run every 33 minutes and managed file deletions. The 'resolvedns.bat' initiated the download of a file from the C2 server. The 'versioninfo.bat' collected information about the user's computer and exfiltrated data using 'unactivate.vbs'.",
        "Persistence": "The 'stopedge.bat' file contributed to persistence by scheduling 'install.vbs' to run every 33 minutes if 'paycom.ini' was present in the directory.",
        "Privilege Escalation": "None",
        "Defense Evasion": "The scripts used in the attack, such as 'install.vbs' and 'update.vbs', were heavily obfuscated using XOR and various techniques to evade detection. The 'stopedge.bat' file also managed the deletion of 'versioninfo.bat' and 'update.vbs' after execution, which is indicative of defense evasion through indicator removal.",
        "Credential Access": "None",
        "Discovery": "The 'versioninfo.bat' script was responsible for discovering information about the user's computer, including directory listings, IP addresses, running processes, and system details.",
        "Lateral Movement": "None",
        "Collection": "The 'versioninfo.bat' script collected a variety of information from the user's computer and prepared it for exfiltration.",
        "Command and Control": "The 'activate.vbs' script communicated with a C2 server by sending an HTTP request and saving the content retrieved. The 'resolvedns.bat' also facilitated communication with the C2 server to initiate the download of a file.",
        "Exfiltration": "The 'unactivate.vbs' script exfiltrated data through a POST request, and the 'versioninfo.bat' script exfiltrated collected information about the user's computer using 'unactivate.vbs'.",
        "Impact": "None",
        "Others": "The Konni APT Group is identified as an Advanced Persistent Threat with origins believed to be in North Korea, primarily targeting South Korea. Their modus operandi typically involves spear-phishing emails and malicious documents. The group's objectives include data exfiltration and espionage. The report includes a detailed breakdown of MITRE ATT&CK techniques employed by the Konni APT Group, as well as mitigations and detection strategies, including a YARA rule for detection."
    },
    "technique_label": true
}